# target/linux/au1000/base-files/etc/uci-defaults/02_network
#!/bin/sh

. /lib/au1000.sh
. /lib/functions/system.sh
. /lib/functions/uci-defaults.sh
MAC=

case "$(au1000_board_name)" in
	'internetbox')
		MAC="$( au1000_yamonenv_getvar 'ethaddr' )"
		MAC="$( macaddr_canonicalize "$MAC" )"

		[ -n "$MAC ] && ucidef_set_interface_macaddr lan "$MAC"
	;;
esac

[ -z "$MAC" ] || uci commit network
