# target/linux/ar71xx/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/ar71xx.sh

board=$(ar71xx_board_name)

case "$board" in
alfa-nx)
	ucidef_set_led_netdev "wan" "WAN" "alfa:green:led_2" "eth0"
	ucidef_set_led_netdev "lan" "LAN" "alfa:green:led_3" "eth1"
	;;

all0258n)
	ucidef_set_rssimon "wlan0" "40000" "1"
	ucidef_set_led_rssi "rssilow" "RSSILOW" "all0258n:red:rssilow" "wlan0" "1" "40" "0" "6"
	ucidef_set_led_rssi "rssimedium" "RSSIMEDIUM" "all0258n:yellow:rssimedium" "wlan0" "30" "80" "-29" "5"
	ucidef_set_led_rssi "rssihigh" "RSSIHIGH" "all0258n:green:rssihigh" "wlan0" "70" "100" "-69" "8"
	;;

all0315n)
	ucidef_set_rssimon "wlan0" "40000" "1"
	ucidef_set_led_rssi "rssilow" "RSSILOW" "all0315n:red:rssilow" "wlan0" "1" "40" "0" "6"
	ucidef_set_led_rssi "rssimedium" "RSSIMEDIUM" "all0315n:yellow:rssimedium" "wlan0" "30" "80" "-29" "5"
	ucidef_set_led_rssi "rssihigh" "RSSIHIGH" "all0315n:green:rssihigh" "wlan0" "70" "100" "-69" "8"
	;;

ap113)
	ucidef_set_led_usbdev "usb" "USB" "ap113:green:usb" "1-1"
	;;

cap4200ag)
	ucidef_set_led_default "lan_green" "LAN_GREEN" "senao:green:lan" "1"
	ucidef_set_led_wlan "wlan_amber" "WLAN_AMBER" "senao:amber:wlan" "phy0tpt"
	ucidef_set_led_wlan "wlan_green" "WLAN_GREEN" "senao:green:wlan" "phy1tpt"
	;;

db120)
	ucidef_set_led_usbdev "usb" "USB" "db120:green:usb" "1-1"
	;;

rb-750)
	ucidef_set_led_default "act" "act" "rb750:green:act" "1"
	ucidef_set_led_netdev "port1" "port1" "rb750:green:port1" "eth1"
	ucidef_set_led_switch "port2" "port2" "rb750:green:port2" "switch0" "0x10"
	ucidef_set_led_switch "port3" "port3" "rb750:green:port3" "switch0" "0x08"
	ucidef_set_led_switch "port4" "port4" "rb750:green:port4" "switch0" "0x04"
	ucidef_set_led_switch "port5" "port5" "rb750:green:port5" "switch0" "0x02"
	;;

dir-600-a1|\
dir-615-e4)
	ucidef_set_led_netdev "wan" "WAN" "d-link:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "d-link:green:lan1" "switch0" "0x02"
	ucidef_set_led_switch "lan2" "LAN2" "d-link:green:lan2" "switch0" "0x04"
	ucidef_set_led_switch "lan3" "LAN3" "d-link:green:lan3" "switch0" "0x08"
	ucidef_set_led_switch "lan4" "LAN4" "d-link:green:lan4" "switch0" "0x10"
	;;

dir-825-b1)
	ucidef_set_led_usbdev "usb" "USB" "d-link:blue:usb" "1-1"
	;;

dir-825-c1)
	ucidef_set_led_usbdev "usb" "USB" "d-link:blue:usb" "1-1"
	;;

hornet-ub)
	ucidef_set_led_netdev "lan" "LAN" "alfa:blue:lan" "eth0"
	ucidef_set_led_netdev "wan" "WAN" "alfa:blue:wan" "eth1"
	ucidef_set_led_wlan "wlan" "WLAN" "alfa:blue:wlan" "phy0tpt"
	ucidef_set_led_usbdev "usb" "USB" "alfa:blue:usb" "1-1"
	;;

mr600)
	ucidef_set_led_wlan "wlan58" "WLAN58" "mr600:green:wlan58" "phy0tpt"
	;;

mzk-w04u)
	ucidef_set_led_usbdev "usb" "USB" "planex:green:usb" "1-1"
	;;

mzk-w300nh)
	ucidef_set_led_wlan "wlan" "WLAN" "planex:green:wlan" "phy0tpt"
	;;

nbg460n_550n_550nh)
	ucidef_set_led_wlan "wlan" "WLAN" "nbg460n:green:wlan" "phy0tpt"
	;;

om2p | \
om2p-hs | \
om2p-lc)
	ucidef_set_led_netdev "port1" "port1" "om2p:blue:wan" "eth0"
	ucidef_set_led_netdev "port2" "port2" "om2p:blue:lan" "eth1"
	;;

tew-712br)
	ucidef_set_led_netdev "wan" "WAN" "trendnet:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "trendnet:green:lan1" "switch0" "0x02"
	ucidef_set_led_switch "lan2" "LAN2" "trendnet:green:lan2" "switch0" "0x04"
	ucidef_set_led_switch "lan3" "LAN3" "trendnet:green:lan3" "switch0" "0x08"
	ucidef_set_led_switch "lan4" "LAN4" "trendnet:green:lan4" "switch0" "0x10"
	ucidef_set_led_wlan "wlan" "WLAN" "trendnet:green:wlan" "phy0tpt"
	;;

tl-mr11u | \
tl-mr3020 | \
tl-mr3040)
	ucidef_set_led_usbdev "usb" "USB" "tp-link:green:3g" "1-1"
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	ucidef_set_led_netdev "lan" "LAN" "tp-link:green:lan"  "eth0"
	;;

tl-mr3220 | \
tl-mr3420 )
	ucidef_set_led_usbdev "usb" "USB" "tp-link:green:3g" "1-1"
	;;

tl-mr3220-v2)
	ucidef_set_led_netdev "wan" "WAN" "tp-link:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "tp-link:green:lan1" "switch0" "0x04"
	ucidef_set_led_switch "lan2" "LAN2" "tp-link:green:lan2" "switch0" "0x08"
	ucidef_set_led_switch "lan3" "LAN3" "tp-link:green:lan3" "switch0" "0x10"
	ucidef_set_led_switch "lan4" "LAN4" "tp-link:green:lan4" "switch0" "0x02"
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	ucidef_set_led_usbdev "usb" "USB" "tp-link:green:3g" "1-1"
	;;

tl-wa901nd)
	ucidef_set_led_netdev "lan" "LAN" "tp-link:green:lan" "eth0"
	;;

tl-wa901nd-v2)
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	;;

tl-wdr4300)
	ucidef_set_led_usbdev "usb1" "USB1" "tp-link:green:usb1" "1-1.1"
	ucidef_set_led_usbdev "usb2" "USB2" "tp-link:green:usb2" "1-1.2"
	ucidef_set_led_wlan "wlan2g" "WLAN2G" "tp-link:blue:wlan2g" "phy0tpt"
	;;

tl-wr741nd)
	ucidef_set_led_netdev "wan" "WAN" "tp-link:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "tp-link:green:lan1" "switch0" "0x02"
	ucidef_set_led_switch "lan2" "LAN2" "tp-link:green:lan2" "switch0" "0x04"
	ucidef_set_led_switch "lan3" "LAN3" "tp-link:green:lan3" "switch0" "0x08"
	ucidef_set_led_switch "lan4" "LAN4" "tp-link:green:lan4" "switch0" "0x10"
	;;

tl-wr741nd-v4)
	ucidef_set_led_netdev "wan" "WAN" "tp-link:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "tp-link:green:lan1" "switch0" "0x04"
	ucidef_set_led_switch "lan2" "LAN2" "tp-link:green:lan2" "switch0" "0x08"
	ucidef_set_led_switch "lan3" "LAN3" "tp-link:green:lan3" "switch0" "0x10"
	ucidef_set_led_switch "lan4" "LAN4" "tp-link:green:lan4" "switch0" "0x02"
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	;;

tl-wr841n-v8)
	ucidef_set_led_netdev "wan" "WAN" "tp-link:green:wan" "eth0"
	ucidef_set_led_switch "lan1" "LAN1" "tp-link:green:lan1" "switch0" "0x04"
	ucidef_set_led_switch "lan2" "LAN2" "tp-link:green:lan2" "switch0" "0x08"
	ucidef_set_led_switch "lan3" "LAN3" "tp-link:green:lan3" "switch0" "0x10"
	ucidef_set_led_switch "lan4" "LAN4" "tp-link:green:lan4" "switch0" "0x02"
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	;;

tl-wr941nd | \
tl-wr1041n-v2)
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	;;

tl-wr1043nd)
	ucidef_set_led_usbdev "usb" "USB" "tp-link:green:usb" "1-1"
	ucidef_set_led_wlan "wlan" "WLAN" "tp-link:green:wlan" "phy0tpt"
	;;

tl-wr2543n)
	ucidef_set_led_usbdev "usb" "USB" "tp-link:green:usb" "1-1"
	;;

wrt160nl)
	ucidef_set_led_wlan "wlan" "WLAN" "wrt160nl:blue:wlan" "phy0tpt"
	;;

wndr3700)
	ucidef_set_led_default "wan" "WAN LED (green)" "wndr3700:green:wan" "0"
	ucidef_set_led_usbdev "usb" "USB" "wndr3700:green:usb" "1-1"
	;;

wndr4300)
	ucidef_set_led_netdev "wan" "WAN (green)" "netgear:green:wan" "eth0.2"
	ucidef_set_led_usbdev "usb" "USB" "netgear:green:usb" "1-1"
	;;

whr-g301n |\
whr-hp-g300n |\
whr-hp-gn)
	ucidef_set_led_netdev "wan" "WAN" "buffalo:green:wan" "eth1"
	ucidef_set_led_switch "lan1" "LAN1" "buffalo:green:lan1" "switch0" "0x02"
	ucidef_set_led_switch "lan2" "LAN2" "buffalo:green:lan2" "switch0" "0x04"
	ucidef_set_led_switch "lan3" "LAN3" "buffalo:green:lan3" "switch0" "0x08"
	ucidef_set_led_switch "lan4" "LAN4" "buffalo:green:lan4" "switch0" "0x10"
	;;

wlae-ag300n)
	ucidef_set_led_netdev "wireless" "WIRELESS" "buffalo:green:wireless" "wlan0"
	;;

wzr-hp-ag300h)
	ucidef_set_led_default "diag" "DIAG" "buffalo:red:diag" "0"
	ucidef_set_led_netdev "router" "ROUTER" "buffalo:green:router" "eth1"
	ucidef_set_led_usbdev "usb" "USB" "buffalo:green:usb" "1-1"
	;;

wzr-hp-g300nh)
	ucidef_set_led_wlan "wlan" "Wireless" "buffalo:green:wireless" "phy0tpt"
	ucidef_set_led_netdev "router" "Router" "buffalo:green:router" "eth1"
	ucidef_set_led_usbdev "usb" "USB" "buffalo:blue:usb" "1-1"
	;;

zcn-1523h-2)
	ucidef_set_led_netdev "lan1" "lan1" "zcn-1523h:green:lan1" "eth0"
	;;

zcn-1523h-5)
	ucidef_set_led_netdev "lan1" "lan1" "zcn-1523h:green:lan1" "eth0"
	ucidef_set_led_netdev "lan2" "lan2" "zcn-1523h:green:lan2" "eth1"
	;;
esac

ucidef_commit_leds

exit 0
