# target/linux/ar71xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/ar71xx.sh

ucidef_set_interface_loopback

board=$(ar71xx_board_name)

case "$board" in
all0315n |\
all0258n |\
ja76pf2|\
ubnt-unifi-outdoor)
	ucidef_set_interface_lan "eth0 eth1"
	;;

ap132)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 2 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 1"
	;;

ap136-010)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "5 6"
	;;

ap136-020 |\
ap135-020)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "5 6"
	;;

db120 |\
rb-2011l | \
rb-2011uas-2hnd)
	ucidef_set_interfaces_lan_wan "eth0.1 eth1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 2 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 1"
	ucidef_add_switch "switch1" "1" "1"
	ucidef_add_switch_vlan "switch1" "1" "0 1 2 3 4 5"
	;;

dir-825-b1|\
tew-673gru)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5t"
	;;

nbg460n_550n_550nh)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "0" "0 1 2 3 5"
	;;

rb-433 |\
rb-433u)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "1 2 5"
	;;

rb-435g)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2"
	;;

rb-450)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
	;;

rb-450g |\
routerstation-pro)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	;;

ap136 |\
rb-750gl |\
rb-751g |\
rb-951g-2hnd |\
wzr-hp-g450h)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 2 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 1"
	;;

rb-493g)
	ucidef_set_interfaces_lan_wan "eth0 eth1.1" "eth1.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	ucidef_add_switch "switch1" "1" "1"
	ucidef_add_switch_vlan "switch1" "1" "0t 1 2 3 4"
	ucidef_add_switch_vlan "switch1" "2" "0t 5"
	;;

wzr-hp-g300nh2 |\
pb92 |\
ap113)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 1 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 2"
	;;

tl-wdr4300|\
tl-wr1041n-v2)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 2 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 1"
	;;

tl-wr1043nd)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "1 2 3 4 5t"
	ucidef_add_switch_vlan "switch0" "2" "0 5t"
	;;

tl-wr2543n)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "1 2 3 4 9t"
	ucidef_add_switch_vlan "switch0" "2" "0 9t"
	;;

tl-wr841n-v1|\
tl-wr941nd)
	ucidef_set_interface_raw "eth" "eth0"
	ucidef_set_interfaces_lan_wan "lan1 lan2 lan3 lan4" "wan"
	;;

tl-mr3420-v2 |\
tl-wr841n-v8)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	;;

uap-pro)
	ucidef_set_interface_lan "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2"
	;;

wrt160nl)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4 5"
	;;

wzr-hp-g300nh)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5t"
	;;

dir-825-c1 |\
dir-835-a1 |\
wndr4300)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "0t 5"
	;;

all0305 |\
aw-nr580 |\
bullet-m |\
cap4200ag |\
eap7660d |\
mr600 |\
mr600v2 |\
rb-411 |\
tl-mr11u |\
tl-mr3020 |\
tl-mr3040 |\
tl-wa901nd |\
tl-wa901nd-v2 |\
tl-wr703n |\
wp543)
	ucidef_set_interface_lan "eth0"
	;;

alfa-ap96 |\
alfa-nx |\
ap83 |\
jwap003 |\
pb42 |\
pb44 |\
routerstation|\
wpe72)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;

ap121 |\
ap121-mini |\
ap96 |\
airrouter |\
dir-600-a1 |\
dir-615-c1 |\
dir-615-e4 |\
ja76pf |\
rb-750 |\
rb-751 |\
tew-632brp |\
tew-712br |\
tl-mr3220 |\
tl-mr3220-v2 |\
tl-mr3420 |\
tl-wdr3500 |\
tl-wr741nd |\
tl-wr741nd-v4 |\
tl-wr841n-v7 |\
whr-g301n |\
whr-hp-g300n |\
whr-hp-gn |\
wzr-hp-ag300h)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 4"
	;;

wzr-hp-g450h)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0t 2 3 4 5"
	ucidef_add_switch_vlan "switch0" "2" "0t 1"
	;;

ew-dorin)
	ucidef_set_interface_lan "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3"
	;;

ew-dorin-router)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 2 3"
	;;

*)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	;;
esac

uci commit network

exit 0
