# target/linux/ar71xx/base-files/etc/uci-defaults/99_zsun-enable-wifi
#!/bin/sh

uci set wireless.@wifi-device[0].disabled=0
uci commit wireless
wifi

exit 0