# target/linux/ar71xx/base-files/etc/uci-defaults/nbg460n_550n_550nh
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

nbg460n_550n_550nh_set_wlan_led() {
	uci batch <<EOF
set system.wlan_led=led
set system.wlan_led.name='WLAN'
set system.wlan_led.sysfs='nbg460n:green:wlan'
set system.wlan_led.trigger='phy0tpt'
commit system
EOF
}

if [ "${board}" == "nbg460n_550n_550nh" ]; then
	nbg460n_550n_550nh_set_wlan_led
fi
