# target/linux/ar71xx/base-files/etc/uci-defaults/rb750
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

rb750_set_leds() {
	uci batch <<EOF
set system.led_act='led'
set system.led_act.name='act'
set system.led_act.sysfs='rb750:green:act'
set system.led_act.default='1'
set system.led_port1='led'
set system.led_port1.name='port1'
set system.led_port1.sysfs='rb750:green:port1'
set system.led_port1.trigger='netdev'
set system.led_port1.dev=eth1
set system.led_port1.mode='tx rx'
set system.led_port2='led'
set system.led_port2.name='port3'
set system.led_port2.sysfs='rb750:green:port2'
set system.led_port2.trigger='netdev'
set system.led_port2.dev=port2
set system.led_port2.mode='link tx rx'
set system.led_port3='led'
set system.led_port3.name='port3'
set system.led_port3.sysfs='rb750:green:port3'
set system.led_port3.trigger='netdev'
set system.led_port3.dev=port3
set system.led_port3.mode='link tx rx'
set system.led_port4='led'
set system.led_port4.name='port4'
set system.led_port4.sysfs='rb750:green:port4'
set system.led_port4.trigger='netdev'
set system.led_port4.dev=port4
set system.led_port4.mode='link tx rx'
set system.led_port5='led'
set system.led_port5.name='port5'
set system.led_port5.sysfs='rb750:green:port5'
set system.led_port5.trigger='netdev'
set system.led_port5.dev=port5
set system.led_port5.mode='link tx rx'
commit system
EOF
}

if [ "${board}" == "rb-750" ]; then
	rb750_set_leds
fi
