# target/linux/ar71xx/base-files/etc/uci-defaults/rb750
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

rb750_set_leds() {
	uci batch <<EOF
set system.led_act='led'
set system.led_act.name='act'
set system.led_act.sysfs='rb750:green:act'
set system.led_act.default='1'
set system.led_port1='led'
set system.led_port1.name='port1'
set system.led_port1.sysfs='rb750:green:port1'
set system.led_port1.trigger='netdev'
set system.led_port1.dev=eth1
set system.led_port1.mode='link tx rx'
commit system
EOF
}

if [ "${board}" == "rb-750" ]; then
	rb750_set_leds
fi
