# target/linux/ar71xx/base-files/etc/uci-defaults/ubootenv
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

add_ubootenv() {
	local dev=$1
	local offset=$2
	local envsize=$3
	local secsize=$4
	local numsec=$5
	uci batch <<EOF
add ubootenv ubootenv
set ubootenv.@ubootenv[-1].dev='$dev'
set ubootenv.@ubootenv[-1].offset='$offset'
set ubootenv.@ubootenv[-1].envsize='$envsize'
set ubootenv.@ubootenv[-1].secsize='$secsize'
set ubootenv.@ubootenv[-1].numsec='$numsec'
EOF
}

[ -e /etc/config/ubootenv ] && exit 0

touch /etc/config/ubootenv

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

case "$board" in
all0258n)
	add_ubootenv /dev/mtd1 0x0 0x10000 0x10000
	;;

esac

uci commit ubootenv
