# target/linux/ar71xx/base-files/etc/uci-defaults/wndr3700
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

wndr3700_set_wan_led() {
	uci batch <<EOF
set system.wan_led=led
set system.wan_led.sysfs='wndr3700:green:wan'
set system.wan_led.name='WAN LED (green)'
set system.wan_led.default='0'
set system.usb_led=led
set system.usb_led.name='USB'
set system.usb_led.sysfs='wndr3700:green:usb'
set system.usb_led.trigger='usbdev'
set system.usb_led.dev='1-1'
set system.usb_led.interval='50'
commit system
EOF
}

if [ "${board}" == "wndr3700" -o "${board}" == "wndr3700v2" ]; then
	wndr3700_set_wan_led
fi
