# target/linux/ar71xx/base-files/etc/uci-defaults/zcn-1523h-2
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

zcn_set_leds() {
	uci batch <<EOF
set system.led_lan1='led'
set system.led_lan1.name='lan1'
set system.led_lan1.sysfs='zcn-1523h:green:lan1'
set system.led_lan1.trigger='netdev'
set system.led_lan1.dev='eth0'
set system.led_lan1.mode='link tx rx'
commit system
EOF
}

if [ "${board}" == "zcn-1523h-2" ]; then
	zcn_set_leds
fi
