# package/bluewan-auto/files/etc/hotplug.d/iface/99-bluewan-check
#!/bin/sh
[ "$ACTION" = ifup ] || exit 0

# 维持防火墙绑定
WAN_IDX=$(uci show firewall | grep "=zone" | grep "'wan'" | cut -d'.' -f2 | head -n1)
if [ -n "$WAN_IDX" ]; then
    for iface in BLUE4WAN; do
        if ! uci get firewall.$WAN_IDX.network | grep -q "$iface"; then
            uci add_list firewall.$WAN_IDX.network="$iface"
            echo "[bluewan-check] 已将 $iface 加入 wan 防火墙区"
        fi
    done
    uci commit firewall
    /etc/init.d/firewall reload
fi

# 调用 IPv6 自动配置脚本
/usr/bin/bluewan-ipv6-setup.sh &

exit 0

