# target/linux/bcm53xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh

ucidef_set_interface_loopback

ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
ucidef_add_switch "switch0" "1" "1"
ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5t"
ucidef_add_switch_vlan "switch0" "2" "4 5t"

uci commit network

exit 0
