# target/linux/bcm53xx/base-files/etc/uci-defaults/05_fix-compat-version

uci set system.@system[0].compat_version="1.1"
uci commit system

exit 0
