# target/linux/bcm53xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh

mtd fixtrx firmware || mtd fixseama firmware
