# target/linux/mvebu/cortexa9_neon/base-files/etc/uci-defaults/03_wireless
#
# Copyright (C) 2015 OpenWrt.org
#

[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
linksys,wrt1200ac|\
linksys,wrt1900ac-v2|\
linksys,wrt1900acs|\
linksys,wrt32x)
	SKU=$(strings /dev/mtd3|sed -ne 's/^cert_region=//p')
	WIFIMAC2G=$(macaddr_add $(cat /sys/class/net/eth0/address) +1)
	WIFIMAC5G=$(macaddr_add $WIFIMAC2G +1)
	case "$SKU" in
		AP)
			REGD=CN
		;;
		AU)
			REGD=AU
		;;
		CA)
			REGD=CA
		;;
		EU)
			REGD=DE
		;;
		US)
			REGD=US
		;;
	esac

	uci get wireless.radio0.country || uci set wireless.radio0.country=$REGD
	uci get wireless.@wifi-iface[0].macaddr || uci set wireless.@wifi-iface[0].macaddr=$WIFIMAC5G
	uci get wireless.radio1.country || uci set wireless.radio1.country=$REGD
	uci get wireless.@wifi-iface[1].macaddr || uci set wireless.@wifi-iface[1].macaddr=$WIFIMAC2G
	;;
esac

uci commit wireless

exit 0
