# target/linux/tegra/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
compulab,trimslice)
	migrate_leds "trimslice:green:right=green:indicator-2" "trimslice:green:left=green:indicator-3"
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
