# target/linux/ar71xx/base-files/etc/init.d/defconfig
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

. /lib/ar71xx.sh

START=05

start() {
	local board=$(ar71xx_board_name)

	[ ! -d /etc/config/defaults/$board ] && board="generic"

	for f in $( ls /etc/config/defaults/$board ); do
		if [ ! -e /etc/config/$f ]; then
			cp /etc/config/defaults/$board/$f /etc/config/
		fi
	done
}
