# target/linux/ar71xx/base-files/etc/init.d/defconfig
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=05

start() {
	. /lib/ar71xx.sh

	local board=$(ar71xx_board_name)

	[ ! -d /etc/defconfig/$board ] && return 0

	for f in $( ls /etc/defconfig/$board ); do
		if [ ! -e /etc/config/$f ]; then
			cp /etc/defconfig/$board/$f /etc/config/
		fi
	done
}
