# package/rooter/0optionalapps/bwallocate/files/etc/init.d/textbwint
#!/bin/sh /etc/rc.common

START=88

log() {
	logger -t "TEXTING" "$@"
}

start()
{
	uci set bwmon.general.enabled="1"
	uci commit bwmon
	/usr/lib/bwmon/textbw.sh &
}
