# target/linux/sunxi/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2013-2014 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/sunxi.sh
. /lib/functions/uci-defaults.sh

ucidef_set_interface_loopback

case "$( sunxi_board_name )" in
"olinuxino-micro"*)
	ucidef_set_interface_lan 'wlan0'
	;;
*)
	ucidef_set_interface_lan 'eth0'
	;;
esac

uci commit network

exit 0
