# target/linux/x86/geos/base-files/etc/uci-defaults/ppp
#!/bin/sh

uci batch <<__EOF__

set network.wan=interface
set network.wan.proto='pppoa'
set network.wan.atmdev=0
set network.wan.card=0
set network.wan.vpi=0
set network.wan.vci=32
set network.wan.encaps='vc'
set network.wan.username='username@myisp.net'
set network.wan.password='password'
set network.wan.keepalive=3
set network.wan.defaultroute=1
set network.wan.peerdns=1
set network.wan.pppd_options='kdebug 1 debug noccp'
set network.wan.auto=1

commit network
__EOF__
