# target/linux/x86/geos/base-files/etc/uci-defaults/wireless
#!/bin/sh

[ -f /etc/config/wireless ] || echo "package wireless" > /etc/config/wireless

uci batch <<__EOF__

delete wireless.radio0
delete wireless.@wifi-iface[0]

set wireless.wlan0=wifi-device
set wireless.wlan0.type=mac80211
set wireless.wlan0.phy=phy0
set wireless.wlan0.channel=11
set wireless.wlan0.hwmode=11g
set wireless.wlan0.country=US

add wireless wifi-iface
set wireless.@wifi-iface[-1].device=wlan0
set wireless.@wifi-iface[-1].network=lan
set wireless.@wifi-iface[-1].mode=ap
set wireless.@wifi-iface[-1].ssid=typhoon
set wireless.@wifi-iface[-1].encryption=psk2+aes
set wireless.@wifi-iface[-1].key="passphrase"
set wireless.@wifi-iface[-1].hidden=0

commit wireless
__EOF__
