# target/linux/qualcommax/ipq60xx/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	edgecore,eap101)
		avail=$(fw_printenv -n upgrade_available)
		[ ${avail} -eq 0 ] && fw_setenv upgrade_available 1
		fw_setenv bootcount 0
	;;
	linksys,mr7350|\
	linksys,mr7500)
		mtd resetbc s_env || true
	;;
	yuncore,fap650)
		fw_setenv owrt_bootcount 0
	;;
	netgear,wax610|\
	netgear,wax610y)
		fw_setenv boot_count 0
	;;
	esac
}
