# target/linux/kirkwood/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2012-2013 OpenWrt.org
#

. /lib/functions/uci-defaults.sh

case "`cat /proc/device-tree/model`" in
"Seagate FreeAgent Dockstar")
	ucidef_set_led_default "health" "status:green:health" "1"
	ucidef_set_led_default "fault" "status:orange:fault" "1"
	;;
"Cloud Engines Pogoplug E02")
	ucidef_set_led_default "health" "status:green:health" "1"
	ucidef_set_led_default "fault" "status:orange:fault" "1"
	;;
*)
	;;
esac

ucidef_commit_leds

exit 0
