# target/linux/kirkwood/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2012-2014 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

set_lan_dhcp() {
	local ifname=$1
	uci batch <<EOF
set network.lan='interface'
set network.lan.ifname='$ifname'
set network.lan.proto='dhcp'
set network.lan6='interface'
set network.lan6.ifname='@lan'
set network.lan6.proto='dhcpv6'
set network.lan6.reqprefix='no'
EOF
}

. /lib/functions/uci-defaults.sh
. /lib/kirkwood.sh

board=$(kirkwood_board_name)

ucidef_set_interface_loopback

case "$board" in
"dockstar"|\
"goflexhome"|\
"goflexnet"|\
"iconnect"|\
"ib62x0"|\
"pogo_e02")
	set_lan_dhcp "eth0"
	;;
"ea3500")
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
	ucidef_add_switch_vlan "switch0" "2" "4 6"
	;;
"ea4500")
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
	ucidef_add_switch_vlan "switch0" "2" "4 6"
	;;
"guruplug-server-plus")
	set_lan_dhcp "eth0" "eth1"
	;;
"sheevaplug" | \
"sheevaplug-esata")
	set_lan_dhcp "eth0"
	;;
*)
	ucidef_set_interface_lan "eth0"
	;;
esac

uci commit network

exit 0
