# target/linux/kirkwood/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

set_lan_dhcp() {
	local ifname=$1
	uci batch <<EOF
set network.lan='interface'
set network.lan.ifname='$ifname'
set network.lan.proto='dhcp'
EOF
}

. /lib/functions/uci-defaults.sh
. /lib/kirkwood.sh

touch /etc/config/network

ucidef_set_interface_loopback

hardware=$(kirkwood_hardware_name)

case "$hardware" in
"Seagate FreeAgent DockStar")
	set_lan_dhcp "eth0"
	;;

"RaidSonic ICY BOX IB-NAS6210")
	set_lan_dhcp "eth0"
	;;

*)
	ucidef_set_interface_lan "eth0"
	;;
esac

uci commit network

exit 0
