# target/linux/lantiq/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
# based on ar71xx
#

. /lib/functions/uci-defaults.sh
. /lib/functions/lantiq.sh

[ -e "/sys/class/leds/wifi" ] && ucidef_set_led_wlan "wifi" "wifi" "wifi" "phy0tpt"
[ -e "/sys/class/leds/usb" ] && ucidef_set_led_usbdev "usb" "usb" "usb" "1-1"
[ -e "/sys/class/leds/usb2" ] && ucidef_set_led_usbdev "usb2" "usb2" "usb2" "2-1"
[ -e "/sys/class/leds/internet" ] && ucidef_set_led_netdev "internet" "internet" "internet" "pppoe-wan"

board=$(lantiq_board_name)

case "$board" in
VG3503J)
	ucidef_set_led_netdev "vdsl" "vdsl" "bt:green:dsl" "ptm0"
	;;
BTHOMEHUBV2B)
	ucidef_set_led_default "power" "power" "soc:blue:power" "1"
#	ucidef_set_led_wlan "wifi" "wifi" "soc:blue:wireless" "phy0radio"
	ucidef_set_led_wlan "wifi" "wifi" "soc:blue:wireless" "phy0tpt"
	ucidef_set_led_netdev "internet" "internet" "soc:blue:broadband" "pppoa-wan"
	ucidef_set_led_usbdev "usb" "usb" "soc:blue:phone" "1-1"
	;;
BTHOMEHUBV3A)
	ucidef_set_led_default "power" "power" "soc:blue:power" "1"
	ucidef_set_led_wlan "wifi" "wifi" "soc:blue:wireless" "phy0tpt"
	ucidef_set_led_netdev "internet" "internet" "soc:blue:broadband" "pppoa-wan"
	;;
VGV7510KW22)
	ucidef_set_led_default "power" "power" "power" "1"
	ucidef_set_led_default "power2" "power2" "power2" "0"
	ucidef_set_led_wlan "wifi" "wifi" "wifi" "phy0radio"
	ucidef_set_led_netdev "dsl" "dsl" "dsl" "nas0"
	ucidef_set_led_netdev "internet_green" "internet_green" "internet_green" "pppoe-wan"
	;;
VGV7519)
	ucidef_set_led_default "power" "power" "power" "0"
	ucidef_set_led_default "power2" "power2" "power2" "1"
	ucidef_set_led_wlan "wireless_yellow" "wireless_yellow" "wireless_yellow" "phy0radio"
	ucidef_set_led_netdev "internet_yellow" "internet_yellow" "internet_yellow" "nas0"
	ucidef_set_led_netdev "broadband_yellow" "broadband_yellow" "broadband_yellow" "pppoe-wan"
	;;
P2812HNUF*)
	ucidef_set_led_default "power" "power" "power" "0"
	ucidef_set_led_default "power2" "power2" "power2" "1"
	ucidef_set_led_wlan "wireless_green" "wireless_green" "wireless_green" "phy0radio"
	ucidef_set_led_netdev "dsl" "dsl" "dsl" "nas0"
	ucidef_set_led_netdev "internet_green" "internet_green" "internet_green" "pppoe-wan"
	;;
ARV7519RW22)
	ucidef_set_led_default "power" "power" "power" "1"
	ucidef_set_led_netdev "internet_green" "internet_green" "internet_green" "pppoe-wan"
	ucidef_set_led_netdev "lan" "lan" "lan" "eth0.1"
	;;
*)
	;;
esac

for a in `ls /sys/class/leds/`; do
	grep -q "\[none\]" /sys/class/leds/$a/trigger
	[ $? -eq 0 ] && ucidef_set_led_default $a $a $a `cat /sys/class/leds/$a/brightness`
done

ucidef_commit_leds

exit 0
