# target/linux/lantiq/base-files/etc/uci-defaults/03_wireless-wps
. /lib/functions.sh
. /lib/functions/lantiq.sh

set_wps() {
	uci set wireless.$1.wps_pushbutton=0
	uci set wireless.$1.wps_device_name="OpenWrt - $(lantiq_board_name)"
}

board=$(lantiq_board_id)
case $board in
WBMR)
	config_load wireless
	config_foreach set_wps wifi-iface
	uci commit wireless
	;;
esac
