# target/linux/lantiq/base-files/etc/uci-defaults/leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
# based on ar71xx
#

. /lib/functions/uci-defaults.sh
. /lib/functions/lantiq.sh

[ -e "/sys/class/leds/wifi" ] && ucidef_set_led_wlan "wifi" "wifi" "wifi" "phy0tx"
[ -e "/sys/class/leds/usb" ] && ucidef_set_led_usbdev "usb" "usb" "usb" "1-1"
[ -e "/sys/class/leds/dsl" ] &&	ucidef_set_led_netdev "dsl" "dsl" "dsl" "pppoe-wan"

for a in `ls /sys/class/leds/`; do
	grep -q "\[none\]" /sys/class/leds/$a/trigger
	[ $? -eq 0 ] && ucidef_set_led_default $a $a $a `cat /sys/class/leds/$a/brightness`
done

board=$(lantiq_board_id)

case "$board" in
*)
	;;
esac

ucidef_commit_leds

exit 0
