# target/linux/lantiq/base-files/etc/uci-defaults/leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
# based on ar71xx
#

. /lib/functions/uci-defaults.sh
. /lib/lantiq.sh

board=$(lantiq_board_name)

case "$board" in
ARV4520PW|ARV452CPW)
	ucidef_set_led_netdev "wifi" "wifi" "soc:blue:wifi" "wlan0"
	ucidef_set_led_usbdev "usb" "USB" "soc:blue:usb" "1-1"
	ucidef_set_led_default "fxs1" "fxs1" "soc:blue:fxs1" "1"
	ucidef_set_led_default "fxs2" "fxs2" "soc:blue:fxs2" "1"
	ucidef_set_led_default "online" "online" "soc:blue:internet" "1"
	;;

ARV7525PW)
	ucidef_set_led_netdev "wifi" "wifi" "soc:green:wifi" "wlan0"
	ucidef_set_led_default "fxs1" "fxs1" "soc:green:fxs1" "1"
	ucidef_set_led_default "fxs2" "fxs2" "soc:green:fxs2" "1"
	ucidef_set_led_default "dsl" "dsl" "soc:red:dsl" "1"
	ucidef_set_led_default "online" "online" "soc:green:online" "1"
	;;

ARV4525PW)
	ucidef_set_led_netdev "wifi" "wifi" "soc:green:wifi" "wlan0"
	ucidef_set_led_default "fxs1" "fxs1" "soc:green:fxs1" "1"
	ucidef_set_led_default "fxs2" "fxs2" "soc:green:fxs2" "1"
	ucidef_set_led_default "dsl" "dsl" "soc:green:dsl" "1"
	ucidef_set_led_default "online" "online" "soc:green:online" "1"
	;;

GIGASX76X)
	ucidef_set_led_wlan "wifi" "wifi" "soc:green:wifi" "phy0radio"
	ucidef_set_led_netdev "online" "online" "soc:green:online" "pppoe-wan"
	ucidef_set_led_usbdev "usb" "USB" "soc:green:usb" "1-1"
	;;

ARV4519PW)
	ucidef_set_led_wlan "wifi" "wifi" "soc:green:wifi" "phy0radio"
	ucidef_set_led_netdev "internet" "internet" "soc:green:internet" "pppoe-wan"
	ucidef_set_led_usbdev "usb" "USB" "soc:green:usb" "1-1"
	;;

ARV7518PW)
	ucidef_set_led_wlan "wifi" "wifi" "soc:green:wifi" "phy0radio"
	ucidef_set_led_netdev "internet" "internet" "soc:green:internet" "pppoe-wan"
	ucidef_set_led_usbdev "usb" "USB" "soc:green:usb" "1-1"
	;;
esac

ucidef_commit_leds

exit 0
