# package/rooter/0optionalapps/pingtest/files/etc/init.d/pingtest
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=99
USE_PROCD=1

log() {
	logger -t "pingtest" "$@"
}

start_service()
{
	/usr/lib/custom/test_ping.sh &
}
