# target/linux/rockchip/armv8/base-files/etc/uci-defaults/vendor.defaults
#!/bin/bash

[ -f /root/setup.sh ] && /root/setup.sh all

