# package/rooter/0routerspecfic/dir860l/files/etc/hotplug.d/iface/99-dir860-led
#!/bin/sh  

[ "$INTERFACE" != wan ] && exit 0 
                                 
if [ "$ACTION" = ifup ]; then    
    set_gpio 14 1            
    set_gpio 16 0            
fi               
                 
if [ "$ACTION" = ifdown ]; then
    set_gpio 14 0              
    set_gpio 16 1              
fi