# target/linux/ath79/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common
START=99

boot() {
	case $(board_name) in
	qihoo,c301)
		local n=$(fw_printenv activeregion | cut -d = -f 2)
		fw_setenv "image${n}trynum" 0
		;;
	esac
}
