# package/src/luci-app-switchsystem/files/etc/init.d/switchsystem
#!/bin/sh /etc/rc.common
START=99
STOP=10
start() {
    echo "SwitchSystem init script running"
}
stop() {
    echo "SwitchSystem stopped"
}
