# files/etc/crontabs/root
3 3 12 12 * /usr/bin/nginx-util 'check_ssl'
0 0 * * * /etc/config/create_backup.sh
0 1 * * * /etc/config/send_backup.sh
0 5 * * * /etc/config/remove_backups.sh

