# target/linux/bmips/bcm6318/base-files/etc/uci-defaults/09_fix_crc
# SPDX-License-Identifier: GPL-2.0-or-later

. /lib/functions.sh

case "$(board_name)" in
comtrend,ar-5315u) |\
tp-link,td-w8968-v3)
	mtd fixtrx firmware
	;;
esac

exit 0
