# package/rooter/0optionalapps/ext-menu/files/etc/uci-defaults/63-menu
#!/bin/sh

fv=$(uci -q get custom.menu.full)
if [ $fv = "0" ]; then
	fv="1"
else
	fv="0"
fi

/usr/lib/fullmenu/setmenu.sh $fv
