# target/linux/ramips/base-files/etc/uci-defaults/01_leds
#!/bin/sh

. /lib/functions/uci-defaults.sh
. /lib/ramips.sh

set_usb_led() {
	ucidef_set_led_usbdev "usb" "USB" "$1" "1-1"
}

set_wifi_led() {
	ucidef_set_led_netdev "wifi_led" "wifi" "$1" "wlan0"
}

board=$(ramips_board_name)

case $board in
	3g-6200n)
		set_wifi_led "edimax:amber:wlan"
		set_usb_led "edimax:blue:3g"
		;;
	3g300m)
		set_usb_led "tenda:blue:3g"
		;;
	air3gii)
		set_wifi_led "airlive:green:wlan"
		set_usb_led "airlive:green:mobile"
		;;
	all0256n)
		ucidef_set_rssimon "wlan0" "40000" "1"
		ucidef_set_led_rssi "rssilow" "RSSILOW" "all0256n:green:rssilow" "wlan0" "1" "40" "0" "6"
		ucidef_set_led_rssi "rssimedium" "RSSIMEDIUM" "all0256n:green:rssimed" "wlan0" "30" "80" "-29" "5"
		ucidef_set_led_rssi "rssihigh" "RSSIHIGH" "all0256n:green:rssihigh" "wlan0" "70" "100" "-69" "8"
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	bc2)
		set_usb_led "bc2:blue:usb"
		;;
	br6425)
		set_wifi_led "edimax:orange:wlan"
		;;
	broadway)
		set_usb_led "red:diskmounted"
		set_wifi_led "red:wps_active"
		;;
	br6524n)
		set_wifi_led "edimax:blue:wlan"
		;;
	dir-300-b1|\
	dir-600-b1|\
	dir-600-b2|\
	dir-620-a1)
		# ANDed with vlan switch port 4 led state
		ucidef_set_led_default "wan" "WAN LED (amber)" "d-link:amber:wan" "1"
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	dir-615-d|\
	dir-615-h1)
		ucidef_set_led_default "status" "Status LED (amber)" "d-link:amber:status" "0"
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	esr-9753)
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	f5d8235-v1)
		set_usb_led "f5d8235-v1:blue:storage"
		;;
	f5d8235-v2)
		set_usb_led "f5d8235v2:blue:storage"
		;;
	fonera20n)
		set_usb_led "fonera20n:orange:usb"
		set_wifi_led "fonera20n:orange:wifi"
		;;
	all0239-3g|\
	hw550-3g)
		set_usb_led "hw550-3g:green:usb"
		;;
	mofi3500-3gn)
		set_usb_led "mofi3500-3gn:green:usb"
		;;
	mzk-w300nh2)
		set_wifi_led "mzkw300nh2:amber:wlan"
		;;
	nw718)
		set_usb_led "nw718:amber:usb"
		;;
	omni-emb-hpm)
		ucidef_set_led_default "power" "POWER" "emb:orange:power" "1"
		ucidef_set_led_netdev "eth" "ETH" "emb:green:eth" "eth0"
		set_wifi_led "emb:green:wifi"
		;;
	sl-r7205)
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	v11st-fe)
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	w306r-v20)
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	w502u)
		set_usb_led "alfa:blue:usb"
		set_wifi_led "rt2800pci-phy0::radio"
		;;
	wcr-150gn)
		set_usb_led "wcr150gn:amber:user"
		;;
esac

ucidef_commit_leds

exit 0
