# target/linux/ramips/base-files/etc/uci-defaults/03_wireless
#!/bin/sh
#
# Copyright (C) 2016 OpenWrt.org
#

[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions/system.sh
. /lib/ramips.sh

board=$(ramips_board_name)

case "$board" in
dch-m225)
	uci set wireless.@wifi-device[0].disabled=0
	;;
esac

uci commit wireless

exit 0
