# target/linux/ramips/base-files/etc/uci-defaults/12_failsafe_ap
#!/bin/sh

MAC=$(cat /sys/class/net/eth0/address|awk -F ":" '{print $5""$6}'| tr a-z A-Z)
uci set wireless.failsafe.disabled=0
uci set wireless.failsafe.ssid=OSTAR-$MAC
uci set wireless.failsafe.key=OSTAR@$MAC
uci commit wireless

/etc/init.d/network reload

exit 0
