# target/linux/ramips/base-files/etc/uci-defaults/20_travelmate_enabled
#!/bin/sh
uci batch <<-EOF
	set travelmate.global.trm_enabled='1'
	set travelmate.global.trm_iface='trm_wwan'
	commit travelmate
EOF

/etc/init.d/travelmate restart

exit 0
