# target/linux/ramips/base-files/etc/uci-defaults/30_umdns-enabled
#!/bin/sh
uci batch <<-EOF
	set umdns.@umdns[-1].network='lan trm_wwan wan'
	commit umdns
EOF

/etc/init.d/umdns restart

exit 0
