# target/linux/ramips/base-files/etc/uci-defaults/network
#!/bin/sh
RT3X5X=`cat /proc/cpuinfo | grep RT3.5`
[ -z "${RT3X5X}" ] || {
	uci batch <<EOF
set network.lan.ifname=eth0.1
set network.wan=interface
set network.wan.ifname=eth0.2
set network.wan.proto=dhcp
commit network
EOF
}

uci commit network
