# target/linux/airoha/base-files/etc/uci-defaults/15-w1700k-flow-offload
#!/bin/sh

# Install a hotplug script that reloads the firewall after WiFi
# interfaces come up, so fw4 can discover them and add them to the
# flowtable for hardware bridge offload via PPE.
#
# fw4 auto-discovers bridge member devices for the flowtable, but
# at boot the firewall starts before hostapd creates WiFi interfaces.

mkdir -p /etc/hotplug.d/net
cat > /etc/hotplug.d/net/50-flow-offload-wifi <<'HOTPLUG'
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0

DEV="${DEVICENAME:-$INTERFACE}"
[ -n "$DEV" ] || exit 0

case "$DEV" in
	phy*-ap*) ;;
	*) exit 0 ;;
esac

# Skip if the device is already in the flowtable (fw4 picked it up).
nft list flowtable inet fw4 ft 2>/dev/null | grep -qw "\"$DEV\"" && exit 0

# fw4 auto-discovers bridge members on reload. Serialize via flock so
# concurrent hotplug events don't trigger multiple reloads - the last
# one after all WiFi interfaces are up will be enough.
(
	exec 9>/var/lock/flow-offload-wifi.lock
	flock -n 9 || exit 0
	# Brief settle time so sibling AP interfaces also get included.
	sleep 2
	/etc/init.d/firewall reload >/dev/null 2>&1
) &

exit 0
HOTPLUG
chmod +x /etc/hotplug.d/net/50-flow-offload-wifi

exit 0
