# target/linux/ipq806x/base-files/etc/uci-defaults/network
#!/bin/sh
#
# Copyright (c) 2015 The Linux Foundation. All rights reserved.
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/ipq806x.sh

ucidef_set_interface_loopback

board=$(ipq806x_board_name)

case "$board" in
ap148)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "6 1 2 3 4"
	ucidef_add_switch_vlan "switch0" "2" "0 5"
	;;
*)
	echo "Unsupported hardware. Network interfaces not intialized"
	;;
esac

uci commit network

exit 0
