# target/linux/ath79/nand/base-files/etc/hotplug.d/ieee80211/10-fix-wifi-mac
#!/bin/sh

[ "$ACTION" = add ] || exit

PHYNBR=${DEVPATH##*/phy}

. /lib/functions/system.sh

board=$(board_name)

case $board in
zyxel,nbg6716)
	# Set mac address for ath9k device
	[ "$PHYNBR" -eq 1 ] && \
		macaddr_add "$(mtd_get_mac_ascii u-boot-env ethaddr)" 2 > /sys${DEVPATH}/macaddress
	;;
esac
