# target/linux/x86/geos/base-files/etc/hotplug.d/button/50-reboot
#!/bin/sh

logger -t button -p daemon.info "$BUTTON/$ACTION"

case "$BUTTON/$ACTION" in

reset/released)
  reboot -f
  ;;

esac

exit 0
