# package/base-files/files/etc/uci-defaults/15_migrate-time-zonename
zonename="$(uci -q get system.@system[0].zonename)"
case "$zonename" in
        *[[:space:]]*) uci set system.@system[0].zonename="${zonename// /_}" ;;
esac

exit 0
