# package/base-files/files/etc/uci-defaults/94_switch_ports_status
#!/bin/sh

#disable switch_ports_status by default
/etc/init.d/switch_ports_status disable

exit 0
