# target/linux/ath79/nand/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
meraki,z1)
	migrate_leds ":tricolor0=:tricolor"
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
