# package/shenzt/luci-app-adguardhome-1/root/etc/uci-defaults/40_luci-AdGuardHome
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@AdGuardHome[-1]
	add ucitrack AdGuardHome
	set ucitrack.@AdGuardHome[-1].init=AdGuardHome
	commit ucitrack
  delete firewall.AdGuardHome
	set firewall.AdGuardHome=include
	set firewall.AdGuardHome.type=script
	set firewall.AdGuardHome.path=/usr/share/AdGuardHome/firewall.start
	set firewall.AdGuardHome.reload=1
	commit firewall
EOF

rm -f /tmp/luci-indexcache
exit 0
