# package/madwifi/files/madwifi.init
#!/bin/sh /etc/rc.common
START=20
boot() {
	iwpriv ath0 mode 3
}
