# target/linux/ath79/tiny/base-files/etc/hotplug.d/firmware/10-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath9k-eeprom-pci-0000:00:00.0.bin")
	case $board in
	ubnt,airrouter|\
	ubnt,bullet-m-ar7240|\
	ubnt,bullet-m-ar7241|\
	ubnt,nanobridge-m|\
	ubnt,nanostation-loco-m|\
	ubnt,nanostation-m|\
	ubnt,picostation-m)
		caldata_extract "art" 0x1000 0x1000
		;;
	*)
		caldata_die "board $board is not supported yet"
		;;
	esac
	;;
esac
