# package/base-files/files/etc/hotplug.d/net/00-sysctl
#!/bin/sh

if [ -f /etc/sysctl.conf ] && [ "$ACTION" = add ]; then
	sed -ne "/^[[:space:]]*net\..*\.$DEVICENAME\./p" /etc/sysctl.conf | \
		sysctl -e -p - | logger -t sysctl
fi
