# package/base-files/files/etc/hotplug.d/net/51-docker
#!/bin/sh

[ "${ACTION}" = "add" ] || exit 0
[ "${DEVICENAME}" = "docker0" ] || exit 0

/usr/bin/setup_docker_firewall.sh start &
