# package/utils/busybox/files/cron
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

USE_PROCD=1
PROG=/usr/sbin/crond

start_service () {
	[ -z "$(ls /etc/crontabs/)" ] && return 1

	loglevel=$(uci_get "system.@system[0].cronloglevel")
	mkdir -p /var/spool/cron
	ln -s /etc/crontabs /var/spool/cron/ 2>/dev/null

	procd_open_instance
	procd_set_param command "$PROG" -f -c /etc/crontabs -l ${loglevel:-5}
	procd_close_instance
}
