# package/utils/busybox/files/sysntpd
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=98

USE_PROCD=1
PROG=/usr/sbin/ntpd
SERVICE_PID_FILE=/var/run/sysntpd.pid

start_service() {
	local peers
	local args="-n"
	local enable_server

	config_load system
	config_get peers ntp server
	config_get_bool enable_server ntp enable_server 0

	[ $enable_server -eq 0 -a -z "$peers" ] && return

	procd_open_instance
	procd_set_param command "$PROG" -n
        [ $enable_server -ne 0 ] && procd_append_param command -l
	[ -n "$peers" ] && {
		local peer
		for peer in $peers; do
			procd_append_param command -p $peer
		done
	}
	procd_set_param respawn
	procd_close_instance
}

service_triggers()
{
	procd_add_reload_trigger "system"
}
