# package/ppp/files/etc/hotplug.d/atm/20-atm-modem
#!/bin/sh

if [ "$ACTION" = "add" ]; then
	include /lib/network
	scan_interfaces

	local found=0
	local ifc
	for ifc in $interfaces; do
		local up
		config_get_bool up "$ifc" up 0

		local auto
		config_get_bool auto "$ifc" auto 1

		local proto
		config_get proto "$ifc" proto

		if [ "$proto" = "pppoa" ] && [ "$up" != 1 ] && [ "$auto" = 1 ]; then
			found=1
			( sleep 1; ifup "$ifc" ) &
		fi
	done

	if [ "$found" != 1 ]; then
		logger "Found no matching interface for DSL device $DEVICENAME"
	fi
fi
