# target/linux/ramips/base-files/etc/init.d/enablemodem
#!/bin/sh /etc/rc.common

START=99

start() {
	. /lib/ramips.sh

	local board=$(ramips_board_name)
	if [ $board = "mr200" ]; then
		adb wait-for-device
		adb shell chmod +x /WEBSERVER/www/cgi-bin/*
		adb shell httpd -h /WEBSERVER/www/ &
		echo "2357 000d" > /sys/bus/usb-serial/drivers/option1/new_id
		sleep 2
		adb kill-server
	fi
}
