# target/linux/omap24xx/base-files/etc/init.d/p54spi-eeprom
#!/bin/sh /etc/rc.common

START=30

start() {
	[ -e /lib/firmware/3826.eeprom ] && return 0

	[ -x /usr/bin/calvaria ] || { echo "p54spi EEPROM: calvaria not found"; return 1; }
	[ -x /usr/bin/cal2p54 ] || { echo "p54spi EEPROM: cal2p54 not found"; return 1; }

	/usr/bin/calvaria -p -n wlan-iq-align -i last /dev/mtdblock1 >/tmp/wlan-iq-align || {
		echo "p54spi EEPROM: Failed to extract wlan-iq-align"
		return 1
	}
	/usr/bin/calvaria -p -n wlan-tx-gen2 -i last /dev/mtdblock1 >/tmp/txgen_file || {
		echo "p54spi EEPROM: Failed to extract wlan-tx-gen2"
		return 1
	}

	/usr/bin/cal2p54 /tmp/txgen_file /tmp/wlan-iq-align >/lib/firmware/3826.eeprom || {
		echo "p54spi EEPROM: Failed to generate EEPROM"
		return 1
	}

	rm -f /tmp/wlan-iq-align /tmp/txgen_file
}
