# openwrt/target/default/target_skeleton/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/sbin/halt
tts/0::wait:/bin/login
