# openwrt/target/default/target_skeleton/etc/inittab
::sysinit:/etc/init.d/rcS 2>&1 | logger -p 6 -t '' -s
::shutdown:/sbin/halt
tts/0::askfirst:/bin/ash --login
