# xnd/tip/firstcontact/files/etc/uci-defaults/zzz-firstcontact
#!/bin/sh
[ -f "/etc/ucentral/redirector.json" ] && /etc/init.d/firstcontact disable
[ -f "/etc/ucentral/redirector.json" ] || /etc/init.d/ucentral disable
