# openwrt/package/bitlbee/ipkg/bitlbee.conffiles
/etc/bitlbee/motd.txt
/etc/bitlbee/bitlbee.conf
/etc/init.d/S55bitlbee
