# target/linux/realtek/base-files/etc/init.d/i2c_poe
#!/bin/sh /etc/rc.common

START=10

mp3924_auto_mode() {
	local i2c_addr=0x21

	# Set power limit to 60 Watts.
	i2cset -y         0 ${i2c_addr?} 0x54 0x00 # PMAX[0]   = 0
	i2cset -y         0 ${i2c_addr?} 0x55 0x4b # PMAX[1-8] = 0x4b
	i2cset -y -m 0x04 0 ${i2c_addr?} 0x0f 0x04 # PMAXEN    = 1

	# Set all ports to AUTO mode.
	i2cset -y         0 ${i2c_addr?} 0x03 0xff # MODE    = AUTO
	i2cset -y -m 0x0f 0 ${i2c_addr?} 0x0c 0x0f # 2EVNTEN = 1
	i2cset -y -m 0x0f 0 ${i2c_addr?} 0x07 0x0f # DISEN   = 1
	i2cset -y         0 ${i2c_addr?} 0x06 0xff # DETEN   = 1, CLSEN = 1	
}

boot() {
	case $(board_name) in
	tplink,sg2008p-v3)
		mp3924_auto_mode
		;;
	esac
}
