# target/linux/realtek/base-files/etc/init.d/jg926a_fan_ctrl
#!/bin/sh /etc/rc.common

START=21
USE_PROCD=1

boot() {
	case $(board_name) in
	
logger -t fan.control "Fan Control Script Started"
procd_open_instance
procd_set_param command "/sbin/jg926a_fan_ctrl.sh"
procd_set_param stdout 0
procd_set_param stderr 0
procd_set_param respawn
procd_close_instance
esac
}
