# target/linux/mpc83xx/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyS0::askfirst:/bin/ash --login
ttyS1::askfirst:/bin/ash --login
