# package/independent/base-settings-xwrt-xf-949/files/factory_auth.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org
# Copyright (C) 2019 X-WRT <dev@x-wrt.com>

START=11

start() {
	mt7628factory.sh check || touch /tmp/factory_auth
}

stop() {
	:
}

restart() {
	:
}

reload() {
	:
}
