# target/linux/sibyte/base-files/etc/inittab
# Copyright (C) 2006-2009 OpenWrt.org

::sysinit:/etc/init.d/rcS
duart0::askfirst:/bin/ash --login
