# target/linux/realview/base-files/etc/uci-defaults/02-network
#!/bin/sh

. /lib/functions/uci-defaults.sh

ucidef_set_interface_loopback
ucidef_set_interface_raw "lan" "eth0" "dhcp"
