# target/linux/ath79/generic/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	jjplus,ja76pf2|\
	meraki,mr12|\
	meraki,mr16|\
	ubnt,routerstation|\
	ubnt,routerstation-pro)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
		;;
esac

exit 0
