# target/linux/ath79/generic/base-files/etc/uci-defaults/95_enable-wifi
#!/bin/sh

[ ! -z "$(uci -q get network.lan.ifname)" ] && exit 0

uci set wireless.@wifi-device[0].disabled=0
uci commit wireless

exit 0
