# feeds/luci/applications/luci-app-simple-adblock/root/etc/uci-defaults/40_luci-simple-adblock
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@simple-adblock[-1]
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
