# package/base-files/default/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/sbin/halt
tts/0::askfirst:/bin/ash --login
tts/1::askfirst:/bin/ash --login
