# package/rooter/0routerspecfic/h721/files/etc/init.d/custom
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=10

start()
{
	/usr/lib/custom/custom.lua &
}
