# xnd/tip/tip-defaults/files/etc/init.d/tip-version
#!/bin/sh /etc/rc.common

START=80

boot() {
	echo $(cat /etc/openwrt_release | grep DISTRIB_TIP= | cut -d\' -f2) > /tmp/ucentral.version
}
