# files/etc/config/uci-defaults/96-thisnode-info-translation
#!/bin/sh
sed -i '\|/thisnode.info/|{p;s|/thisnode.info/|/minodo.info/|/meuno.info/|}' /etc/config/dhcp
