# target/linux/qualcommax/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	edgecore,eap102|\
	yuncore,ax880)
		fw_setenv upgrade_available 0
		# Unset changed flag after sysupgrade complete
		fw_setenv changed
	;;
	linksys,mx4200v1|\
	linksys,mx4200v2|\
	linksys,mx4300|\
	linksys,mx5300|\
	linksys,mx8500)
		mtd resetbc s_env || true
	;;
	esac
}
