# target/linux/mediatek/filogic/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-only

START=99

boot() {
	case $(board_name) in
	dlink,aquila-pro-ai-m30-a1)
		if grep -q bootpart=ubi0 /proc/cmdline; then
			fw_setenv bootpart 0
		else
			fw_setenv bootpart 1
		fi
		;;
	zyxel,ex5700-telenor)
		fw_setenv uboot_bootcount 0
		;;
	esac
}
