# target/linux/ath79/base-files/etc/uci-defaults/04_led_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
engenius,epg5000)
	migrate_leds ":wlan-2g=:wlan2g" ":wlan-5g=:wlan5g"
	;;
tplink,archer-c7-v2)
	migrate_leds ":blue:=:green:"
	;;
tplink,archer-c7-v4)
	migrate_leds "archer-c7-v4:=tp-link:"
	;;
tplink,archer-c7-v5)
	migrate_leds "archer-c7-v5:=tp-link:"
	;;
esac

migrations_apply system

exit 0
