# target/linux/ath79/base-files/etc/uci-defaults/04_led_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)
boardonly="${board##*,}"

case "$board" in
engenius,epg5000)
	migrate_leds ":wlan-2g=:wlan2g" ":wlan-5g=:wlan5g"
	;;
tplink,archer-c7-v2)
	migrate_leds ":blue:=:green:"
	;;
tplink,archer-c7-v4|\
tplink,archer-c7-v5)
	migrate_leds "^$boardonly:=tp-link:"
	;;
tplink,re355-v1)
	migrate_leds "re355:=tp-link:"
	;;
tplink,re450-v1)
	migrate_leds "re450:=tp-link:"
	;;
wd,mynet-n750)
	migrate_leds "wd:=mynet-n750:"
	;;
esac

migrations_apply system

exit 0
