# package/miau/ipkg/miau.conffiles
/etc/default/miau
/etc/init.d/S50miau
/etc/miau/miaurc
