# package/rooter/0routerspecfic/we826/files/etc/init.d/wd-init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=50
USE_PROCD=1

log() {
	logger -t "usb-modeswitch" "$@"
}

start_service()
{
	/usr/lib/custom/watchdog.sh &
}
