# package/qos-scripts/files/etc/hotplug.d/iface/10-qos
#!/bin/sh
[ "$ACTION" = ifup ] && /usr/lib/qos.sh interface "$INTERFACE" | sh
