# package/lean/luci-theme-atmaterial2/root/etc/uci-defaults/30_luci-theme-atmaterial2
#!/bin/sh
uci batch <<-EOF
	set luci.themes.AdvancedTomatoMaterial2=/luci-static/atmaterial2
	commit luci
EOF
exit 0
