# target/linux/lantiq/xrx200_legacy/base-files/etc/hotplug.d/firmware/12-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

case "$FIRMWARE" in
	"ath9k-eeprom-pci-0000:00:0e.0.bin" | \
	"ath9k-eeprom-pci-0000:01:00.0.bin" | \
	"ath9k-eeprom-pci-0000:02:00.0.bin")
		board=$(board_name)

		case "$board" in
			tplink,tdw8970|\
			tplink,tdw8980)
				caldata_extract "boardconfig" 0x21000 0x1000
				;;
			*)
				caldata_die "board $board is not supported yet"
				;;
		esac
		;;
esac
