# target/linux/ramips/mt7621_small/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
[ "$ACTION" == "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
	cudy,re1800-v2)
		hw_mac_addr="$(mtd_get_mac_binary bdinfo 0xde00)"
		[ "$PHYNBR" = "1" ] && \
		macaddr_setbit_la "$(macaddr_add $hw_mac_addr 0x100000)" > /sys${DEVPATH}/macaddress
		;;
esac
