# target/linux/rtl838x/base-files/etc/uci-defaults/01_dnsmasq-loopback
#!/bin/sh

uci batch <<EOF
delete dhcp.@dnsmasq[0].interface
add_list dhcp.@dnsmasq[0].interface=lo
commit dhcp
EOF
