# target/linux/msm89xx/base-files/etc/init.d/mount_firmware
#!/bin/sh /etc/rc.common

# This should run before boot service so the firmware is available
# when the kernel modules that need it are loaded.
START=09

start() {
	/sbin/msm-firmware-loader.sh
}

