# xnd/ucentral/udhcpsnoop/files/etc/init.d/dhcpsnoop
#!/bin/sh /etc/rc.common

START=80

USE_PROCD=1
PROG=/usr/sbin/udhcpsnoop

service_triggers() {
	procd_add_reload_trigger dhcpsnoop
}

start_service() {
	[ "$(uci get dhcpsnoop.@snooping[-1].enable)" -eq 1 ] || return
	procd_open_instance
	procd_set_param command "$PROG"
	procd_set_param respawn
	procd_close_instance
}

reload_service() {
	restart
}
