# target/linux/ps3/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
::ctrlaltdel:/sbin/reboot
tts/0::askfirst:/bin/ash --login
ttyS0::askfirst:/bin/ash --login
tty1::respawn:/usr/sbin/initrun
