# target/linux/ath79/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
#!/bin/ash

[ "$ACTION" == "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh
. /lib/functions/k2t.sh

board=$(board_name)

case "$board" in
	phicomm,k2t)
		# The K2T factory firmware does use LAN mac address as the 2.4G wifi mac address
		[ "$PHYNBR" -eq 1 ] && \
			echo $(k2t_get_mac "lan_mac") > /sys${DEVPATH}/macaddress
		;;
	dlink,dir-825-b1)
		[ "$PHYNBR" -eq 0 ] && \
			mac0=$(dd if=$(find_mtd_part "caldata") bs=1 skip=65440 count=17 2>/dev/null)
			echo $mac0 > /sys${DEVPATH}/macaddress

		[ "$PHYNBR" -eq 1 ] && \
			mac1=$(dd if=$(find_mtd_part "caldata") bs=1 skip=65460 count=17 2>/dev/null)
			echo $(macaddr_add "$mac1" 1) > /sys${DEVPATH}/macaddress
		;;
	*)
		;;
esac

