# target/linux/s3c24xx/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
s3c2410_serial0::askfirst:/bin/ash --login
s3c2410_serial2::askfirst:/bin/ash --login
tty1::askfirst:/bin/ash --login
