# target/linux/ps3/petitboot/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
::ctrlaltdel:/sbin/reboot
::restart:/sbin/init
tty1::respawn:/sbin/initrun
tty2::askfirst:/bin/ash --login
tty3::askfirst:/bin/ash --login
