# openwrt/package/madwifi/files/madwifi.init
#!/bin/sh /etc/rc.common
start() {
	iwpriv ath0 mode 3
}
