# package/network/luci-app-n2n/files/etc/uci-defaults/luci-app-n2n
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@n2n[-1]
	add ucitrack n2n
	set ucitrack.@n2n[-1].init=n2n
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
