# package/system/gpio-cdev/nu801/files/nu801.init
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-or-later

START=15
STOP=94
USE_PROCD=1

start_service() {
	. /lib/functions.sh

	procd_open_instance
	procd_set_param command /usr/sbin/nu801 "$(board_name)"
	procd_set_param respawn 5 1 5
	procd_close_instance
}
