# files/etc/uci-defaults/99-plc_enable
chmod +x /etc/init.d/plc_basic
/etc/init.d/plc_basic enable
sleep 3
/etc/init.d/plc_basic start

exit 0
