# package/netifd/files/etc/init.d/network
#!/bin/sh /etc/rc.common
START=40
STOP=90

start() {
	setup_switch() { return 0; }

	include /lib/network
	setup_switch

	ubus call network reload

	grep -qs config /etc/config/wireless && {
		/sbin/wifi up
	}
}

restart() {
	start
}

stop() {
	/sbin/ifdown -a
}
