# package/rooter/ext-rooter-basic/files/etc/hotplug.d/usb/20-usb_mode
#!/bin/sh

/usr/lib/rooter/modeswitch.sh &