# package/base-files/files/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=10

system_config() {
	local cfg="$1"
	local hostname
	
	config_get hostname "$cfg" hostname
	echo "${hostname:-OpenWrt}" > /proc/sys/kernel/hostname

	config_get conloglevel "$cfg" conloglevel
	[ -n "$conloglevel" ] && dmesg -n "$conloglevel"
	
	config_get timezone "$cfg" timezone 
	[ -z "$timezone" ] && timezone=UTC
	echo "$timezone" > /tmp/TZ

	config_get log_ip "$cfg" log_ip 
	config_get log_size "$cfg" log_size
	config_get log_port "$cfg" log_port
	[ -x /sbin/syslogd ] && syslogd -C${log_size:-16} ${log_ip:+-L -R $log_ip${log_port:+:$log_port}}
	[ -x /sbin/klogd ] && klogd
}

apply_uci_config() {(
	include /lib/config
	uci_apply_defaults
)}

start() {
	[ -f /proc/mounts ] || /sbin/mount_root
	[ -f /proc/jffs2_bbc ] && echo "S" > /proc/jffs2_bbc
	[ -f /proc/net/vlan/config ] && vconfig set_name_type DEV_PLUS_VID_NO_PAD

	apply_uci_config
	config_load system
	config_foreach system_config system
	
	mkdir -p /var/run
	mkdir -p /var/log
	mkdir -p /var/lock
	mkdir -p /var/state
	mkdir -p /tmp/.uci
	chmod 0700 /tmp/.uci
	touch /var/log/wtmp
	touch /var/log/lastlog
	touch /tmp/resolv.conf.auto
	ln -sf /tmp/resolv.conf.auto /tmp/resolv.conf
	[ "$FAILSAFE" = "true" ] && touch /tmp/.failsafe

	killall -q hotplug2
	[ -x /sbin/hotplug2 ] && /sbin/hotplug2 --override --persistent \
			--max-children 1 --no-coldplug >/dev/null 2>&1 &

	# the coldplugging of network interfaces needs to happen later, so we do it manually here
	for iface in $(awk -F: '/:/ {print $1}' /proc/net/dev); do
		/usr/bin/env -i ACTION=add INTERFACE="$iface" /sbin/hotplug-call net
	done

	# create /dev/root if it doesn't exist
	[ -e /dev/root ] || {
		rootdev=$(awk 'BEGIN { RS=" "; FS="="; } $1 == "root" { print $2 }' < /proc/cmdline)
		[ -n "$rootdev" ] && ln -s "$rootdev" /dev/root
	}

	load_modules /etc/modules.d/*
}
