# package/base-files/files/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2024 OpenWrt.org

START=10
STOP=90

uci_apply_defaults() {
	. /lib/functions/system.sh
	local IFS=$'\n'
	for f in /etc/uci-defaults/* ; do
		( . "$f" ) && rm -f "$f"
	done
	[ $(uci changes|wc -l) -eq 0 ] || uci commit
}

boot() {
	[ -f /proc/mounts ] || /sbin/mount_root
	[ -f /proc/jffs2_bbc ] && echo S > /proc/jffs2_bbc

	mkdir -p -m 1777 /var/lock
	mkdir -p /var/log
	touch /var/log/lastlog
	touch /var/log/wtmp
	mkdir -p /var/run
	ln -s /var/run /run
	ln -s /var/lock /run/lock
	mkdir -p /var/state
	mkdir -p /var/tmp
	mkdir -p -m 0700 /tmp/.uci
	mkdir -p /tmp/resolv.conf.d
	touch /tmp/resolv.conf.d/resolv.conf.auto
	ln -s /tmp/resolv.conf.d/resolv.conf.auto /tmp/resolv.conf

	local common=nosuid,nodev,noexec,noatime
	[ -d /sys/kernel/debug ] \
		&& mount -t debugfs -o $common debugfs /sys/kernel/debug
	[ -d /sys/fs/bpf ] \
		&& mount -t bpf -o $common,mode=0700 bpffs /sys/fs/bpf
	[ -d /sys/fs/pstore ] \
		&& mount -t pstore -o $common pstore /sys/fs/pstore

	[ "$FAILSAFE" = "true" ] && touch /tmp/.failsafe

	touch /tmp/.config_pending
	/sbin/kmodloader
	# compat for bcm47xx and mvebu
	[ -f /etc/config/wireless ] || sleep 1
	# allows /etc/board.d/* scripts to use values from uboot environment
	[ -f /etc/uci-defaults/30_uboot-envtools ] \
		&& (. /etc/uci-defaults/30_uboot-envtools) 
	/bin/config_generate
	rm -f /tmp/.config_pending

	[ -x /sbin/wifi ] && /sbin/wifi config
	uci_apply_defaults
	sync
	
	# temporary hack until configd exists
	/sbin/reload_config
}
