# package/base-files/files/etc/init.d/done
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=95
boot() {
	mount_root done
	rm -f /sysupgrade.tgz && sync
	rm -f /tmp/sysupgrade.tar && sync

	# process user commands
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}

	if grep -q "openwrt_rootfs" /proc/cmdline; then
		. /lib/upgrade/common.sh
		if is_gpt_dualboot && export_partdevice_label bootpart openwrt_boot; then
			mkdir -p /tmp/boot
			mount "/dev/$bootpart" -o rw,noatime /tmp/boot
			if grep -q openwrt_rootfs_alt /proc/cmdline; then
				grub-editenv /tmp/boot/boot/grub/grubenv set attempt=0 target_slot=1
			else
				grub-editenv /tmp/boot/boot/grub/grubenv set attempt=0 target_slot=0
			fi
			umount /tmp/boot
		fi
	fi

	# set leds to normal state
	. /etc/diag.sh
	set_state done
}
