# package/base-files/files/etc/init.d/initialize.sh
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/lib/hotspot
}
