# package/base-files/files/etc/init.d/network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=40
boot() {
	setup_switch() { return 0; }

	include /lib/network
	setup_switch
	[ -e /etc/config/wireless ] || \
		/sbin/wifi detect > /etc/config/wireless
	/sbin/wifi up
}

start() {
	ifup -a
	/sbin/wifi up
}

restart() {
	ifup -a
	/sbin/wifi up
}

stop() {
	ifdown -a
}
