# package/base-files/files/etc/init.d/network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=40
STOP=40

config_cb() {
	local cfg_type="$1"

	case "$cfg_type" in
		interface)
			append networks "$2" "$N"
		;;
	esac
}

boot() {
	uci_load network
	setup_switch() { return 0; }

	include /lib/network
	setup_switch
	[ -s /etc/config/wireless ] || \
		/sbin/wifi detect > /etc/config/wireless
	/sbin/wifi up
	for network in $networks; do
		config_get startupcheck $network started
		[ "$startupcheck" = "" ] && ifup $network
	done
}

start() {
	ifup -a
	/sbin/wifi up
}

restart() {
	setup_switch() { return 0; }
	
	include /lib/network
	setup_switch
	ifup -a
	/sbin/wifi up
}

stop() {
	ifdown -a
}
