# package/base-files/files/etc/init.d/network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=40
STOP=90

boot() {
	setup_switch() { return 0; }

	include /lib/network
	setup_switch
	/sbin/wifi detect >> /etc/config/wireless
	grep -qs config /etc/config/wireless && {
		/sbin/wifi up
	} || {
		rm -f /etc/config/wireless
	}

	scan_interfaces

	local ifc
	for ifc in $interfaces; do
		local proto
		config_get proto "$ifc" proto

		local auto
		config_get_bool auto "$ifc" auto 1

		type "coldplug_interface_$proto" >/dev/null && [ "$auto" = 1 ] && \
			coldplug_interface_$proto "$ifc"
	done
}

start() {
	setup_switch() { return 0; }
	
	include /lib/network
	setup_switch
	ifup -a
	/sbin/wifi up
}

stop() {
	ifdown -a
}

restart() {
	start
}
