# package/base-files/files/etc/init.d/sysctl
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=0
start() {
	[ -f /etc/sysctl.conf ] && sysctl -p -e >&-
}
