# package/base-files/files/etc/init.d/sysfixtime
#!/bin/sh /etc/rc.common
# Copyright (C) 2013-2014 OpenWrt.org

START=00
STOP=90

RTC_DEV=/dev/rtc0
HWCLOCK=/sbin/hwclock

newest() {
	#find newest file in /etc
	ls -t /etc | head -n 1
}

boot() {
	#set system to rtc time
	[ -e "$RTC_DEV" ] && [ -e "$HWCLOCK" ] && $HWCLOCK -s -u -f $RTC_DEV
	#get timestamp of file
	local time="$(date -r "/etc/$(newest)" +%s)"
	#check system time against file and update system and rtc time if file is newer
	[ $(date +%s) -lt $time ] && date -s @$time && $HWCLOCK -w -u -f $RTC_DEV
}

start() {
	exit 0
}

stop() {
	#fix any file with a timestamp in the future
	for item in $(find /etc) ; do
		[ $(date +%s) -lt $(date -r $item +%s) ] && touch -mc $item && logger -t sysfixtime "Updated timestamp on $item"
	done
	#set rtc to system time
	[ -e "$RTC_DEV" ] && [ -e "$HWCLOCK" ] && $HWCLOCK -w -u -f $RTC_DEV && logger -t sysfixtime "saved '$(date)' to $RTC_DEV"
	#touch latest file
	local file="/etc/$(newest)"
	touch -mc $file && logger -t sysfixtime "Updated timestamp on $file"
}
