# package/base-files/files/etc/init.d/sysntpd
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=98

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1
SERVICE_PID_FILE=/var/run/sysntpd.pid

start() {
	[ -x $PROG ] || return 1

	local peers

	getpeers() {
		config_get peers "$1" server
	}

	config_load system
	config_foreach getpeers timeserver

	if [ -n "$peers" ]; then
		local peer
		local args="-n"
		for peer in $peers; do
			append args "-p $peer"
		done

		service_start /usr/sbin/ntpd $args
	fi
}

stop() {
	service_stop /usr/sbin/ntpd
}
