# package/base-files/files/etc/init.d/usb
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=39
start() {
	[ -d /proc/bus/usb ] && {
		mount -o noatime -t usbfs none /proc/bus/usb
	}		
}
