# package/base-files/files/etc/init.d/watchdog
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=97
start() {
        [ -c /dev/watchdog ] && watchdog -t 5 /dev/watchdog
}
