# target/linux/siflower/base-files/etc/rc.local
#echo 4 > /sys/devices/platform/10800000.gmac/net/eth0/queues/rx-0/rps_cpus
