# files/etc/init.d/set_mac
#!/bin/sh /etc/rc.common

. /lib/functions/mac.sh

START=90
STOP=15

start()
{
    generate_mac $WIFI_OTP_REG0
    generate_mac $SOFTAP_OTP_REG0
    generate_mac $ETH_OTP_REG0

    local WIFI_MAC=$(cat /etc/config/MAC_$WIFI_OTP_REG0)
    local SOFTAP_MAC=$(cat /etc/config/MAC_$SOFTAP_OTP_REG0)
    local ETH_MAC=$(cat /etc/config/MAC_$ETH_OTP_REG0)

    ifconfig eth0 down
    ifconfig eth0 hw ether $ETH_MAC
    ifconfig eth0 up

    ifconfig wlan0 down
    ifconfig wlan0 hw ether $WIFI_MAC
    ifconfig wlan0 up

    uci set wireless.ap.macaddr=$SOFTAP_MAC
    uci commit wireless
}

stop()
{
    echo ""
}
