# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
rgb_control blink 255 255 255 1000 
exit 0
