# target/linux/mvebu/cortexa9/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	fortinet,fg-30e|\
	fortinet,fg-50e|\
	fortinet,fg-51e|\
	fortinet,fg-52e|\
	fortinet,fwf-50e-2r|\
	fortinet,fwf-51e|\
	linksys,wrt1900ac-v1|\
        linksys,wrt32x)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
		;;
esac

exit 0
