# package/shenzt/luci-app-k3usb/files/root/etc/uci-defaults/k3usb
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@k3usb[-1]
	add ucitrack k3usb
	set ucitrack.@k3usb[-1].init=k3usb
	commit ucitrack
EOF

exit 0
