# target/linux/omap4/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyO2::askfirst:/bin/ash --login
tty1::askfirst:/bin/ash --login
