# target/linux/omap35xx/gumstix/base-files/etc/inittab
# inittab for Gumstix board
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
ttyO2::askfirst:/bin/ash --login
tty1::askfirst:/bin/ash --login
