# package/network/ipv6/6distributed/files/6distributed.init
#!/bin/sh /etc/rc.common

START=15

SERVICE_DAEMONIZE=1
SERVICE_WRITE_PID=1
SERVICE_PID_FILE=/var/run/6distributed.pid
DISTRIBUTED_BIN=/usr/sbin/6distributed

start()
{
	service_start $DISTRIBUTED_BIN
}

stop()
{
	service_stop $DISTRIBUTED_BIN
}
