# target/linux/ifxmips/base-files/etc/hotplug.d/button/00-reset
[ "$ACTION" = "released" -a "$BUTTON" = reset ] && {
	reboot
}
