# target/linux/bmips/bcm6328/base-files/etc/uci-defaults/09_fix_crc
# SPDX-License-Identifier: GPL-2.0-or-later

. /lib/functions.sh

case "$(board_name)" in
arcadyan,ar7516 |\
comtrend,ar-5381u |\
comtrend,ar-5387un |\
nucom,r5010unv2)
	mtd fixtrx firmware
	;;
esac

exit 0
