# package/rooter/0optionalapps/ext-zerotier/files/etc/uci-defaults/64-zerotier
#!/bin/sh /etc/rc.common

id=$(uci -q get custom.zerotier.networkid)
uci delete zerotier.zerotier.join
uci add_list zerotier.zerotier.join=$id
uci commit zerotier
/etc/init.d/zerotier restart
