# package/network/config/wifi-scripts/files/etc/hotplug.d/ieee80211/06-wifi-migrate
#!/bin/sh

MIGRATE_SCRIPT="/usr/share/ucode/wifi/migrate.uc"

[ "${ACTION}" = "add" ] && [ -f "$MIGRATE_SCRIPT" ] && \
	/usr/bin/ucode "$MIGRATE_SCRIPT"
