# target/linux/brcm63xx/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2013-2014 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/brcm63xx.sh

case "$board_name" in
963281T_TEF)
	ucidef_set_led_netdev "lan" "LAN" "A4001N1:green:eth" "eth0"
	ucidef_set_led_usbdev "usb" "USB" "A4001N1:green:3g" "1-1"
	;;
'F@ST2704V2')
	ucidef_set_led_usbdev "usb" "USB" "F@ST2704V2:green:usb" "1-1"
	;;
HW553)
	ucidef_set_led_netdev "lan" "LAN" "HW553:blue:lan" "eth0"
	ucidef_set_led_usbdev "usb1" "USB1" "HW553:red:hspa" "1-1"
	ucidef_set_led_usbdev "usb2" "USB2" "HW553:blue:hspa" "1-2"
	;;
HW556*)
	ucidef_set_led_netdev "lan" "LAN" "HW556:red:dsl" "eth0"
	ucidef_set_led_usbdev "usb" "USB" "HW556:red:hspa" "1-2"
	;;
esac

ucidef_commit_leds

exit 0
