# target/linux/brcm63xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/brcm63xx.sh

ucidef_set_interface_loopback

case "$board_name" in

96348W3 |\
96358-502V |\
CVG834G_E15R3921 |\
RTA770BW |\
RTA770W |\
SPW500V)
	ucidef_set_interface_lan "eth0"
	;;

96338GW |\
96338W |\
96338W2_E7T |\
96348A-122 |\
96348W3 |\
AR1004G |\
CPVA642 |\
CT536_CT5621 |\
CT6373-1 |\
D-4P-W |\
"F@ST2604" |\
RTA1320_16M |\
V2110)
	ucidef_set_interface_lan "eth0"
	ucidef_add_switch "eth0" "1" "1"
	ucidef_add_switch_vlan "eth0" "1" "0 1 2 3 4 5"
	;;

96348GW |\
96348GW-10 |\
96348GW-11 |\
96348GW-A |\
96358VW |\
96358VW2 |\
AGPF-S0 |\
"CPVA502+" |\
DV201AMR |\
DWV-S0 |\
"F@ST2404" |\
MAGIC |\
RTA1025W_16 |\
V2500V_BB)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	ucidef_add_switch "eth1" "1" "1"
	ucidef_add_switch_vlan "eth1" "1" "0 1 2 3 4 5"
	;;

AW4139 |\
AW4339U |\
HW553 |\
HW556*)
	ucidef_set_interface_lan "eth0.1"
	ucidef_add_switch "eth0" "1" "1"
	ucidef_add_switch_vlan "eth0" "1" "0 1 2 3 4 5t"
	;;

NB4-* |\
DVG3810BN)
	ucidef_set_interfaces_lan_wan "eth1.1" "eth0"
	ucidef_add_switch "eth1" "1" "1"
	ucidef_add_switch_vlan "eth1" "1" "0 1 2 3 4 5t"
	;;

96328avng |\
96328A-1241N |\
96328A-1441N1 |\
963281TAN |\
963281T_TEF |\
96368M-1341N |\
96368M-1541N |\
96368MVNgr |\
"F@ST2504n" |\
"F@ST2704V2" |\
P870HW-51a_v2)
	ucidef_set_interface_lan "eth0.1"
	ucidef_add_switch "eth0" "1" "1"
	ucidef_add_switch_vlan "eth0" "1" "0 1 2 3 8t"
	;;

96368MVWG)
	ucidef_set_interface_lan "eth0.1"
	ucidef_add_switch "eth0" "1" "1"
	ucidef_add_switch_vlan "eth0" "1" "1 2 4 5 8t"
	;;

96369R-1231N)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "eth0" "1" "1"
	ucidef_add_switch_vlan "eth0" "1" "0 1 2 3 5t"
	ucidef_add_switch_vlan "eth0" "2" "4 5t"
	;;

NB6)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "1 2 3 4 9t"
	ucidef_add_switch_vlan "switch0" "2" "0 9t"
	;;

VW6339GU)
	ucidef_set_interfaces_lan_wan "eth0.1" "eth0.2"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 8t"
	ucidef_add_switch_vlan "switch0" "2" "4 8t"
	;;

963268BU_P300)
	ucidef_set_interface_lan "eth0.1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 3 4 5 6 7 8t"
	;;

*)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;

esac

uci commit network

exit 0
