# target/linux/brcm63xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh
#
# Copyright (C) 2007 OpenWrt.org
#
#

. /lib/brcm63xx.sh

do_fixcrc() {
	mtd fixtrx linux
}

brcm63xx_detect

case "$board_name" in
	96328avng |\
	96328A-1241N |\
	96328A-1441N1 |\
	963281TAN |\
	963281T_TEF |\
	96358-502V |\
	96369R-1231N |\
	AW4339U |\
	"CPVA502+" |\
	CPVA642 |\
	CT6373-1 |\
	MAGIC |\
	RTA770BW |\
	RTA770W |\
	V2110 |\
	V2500V_BB)
		do_fixcrc
		;;
esac

