# target/linux/brcm63xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh
#
# Copyright (C) 2007 OpenWrt.org
#
#

. /lib/brcm63xx.sh

do_fixcrc() {
	mtd fixtrx linux
}

case "$(brcm63xx_board_name)" in
	a4001n |\
	a4001n1 |\
	ar-5381u |\
	ar-5387un |\
	bcm96328avng |\
	bcm963281tan |\
	cpva502p |\
	cpva642 |\
	ct-6373 |\
	dsl-274xb-f |\
	magic |\
	p870hw-51a_v2 |\
	rta770bw |\
	rta770w |\
	spw303v |\
	v2110 |\
	v2500v_bb |\
	vr-3025u |\
	vr-3025un |\
	wap-5813n)
		do_fixcrc
		;;
esac

