# package/mountd/files/mountd.init
#!/bin/sh /etc/rc.common
START=80

start()
{
	/sbin/mountd
}

stop()
{
	killall mountd
}
