# target/linux/cns21xx/base-files/etc/uci-defaults/ns-k330
#!/bin/sh
#
# Copyright (C) 2010 OpenWrt.org
#

. /lib/cns21xx.sh

board=$(get_board_name)

ns_k330_set_leds() {
	uci batch <<EOF
set system.led_eth_orange='led'
set system.led_eth_orange.name='eth_orange'
set system.led_eth_orange.sysfs='ns-k330:orange:eth'
set system.led_eth_orange.trigger='netdev'
set system.led_eth_orange.dev='eth0'
set system.led_eth_orange.mode='link'
set system.led_eth_green='led'
set system.led_eth_green.name='eth_green'
set system.led_eth_green.sysfs='ns-k330:green:eth'
set system.led_eth_green.trigger='netdev'
set system.led_eth_green.dev='eth0'
set system.led_eth_green.mode='tx rx'
set system.usb_led1=led
set system.usb_led1.name='USB'
set system.usb_led1.sysfs='ns-k330:green:usb1'
set system.usb_led1.trigger='usbdev'
set system.usb_led1.dev='1-1'
set system.usb_led1.interval='50'
set system.usb_led2=led
set system.usb_led2.name='USB'
set system.usb_led2.sysfs='ns-k330:green:usb2'
set system.usb_led2.trigger='usbdev'
set system.usb_led2.dev='1-2'
set system.usb_led2.interval='50'
commit system
EOF
}

if [ "${board}" == "ns-k330" ]; then
	ns_k330_set_leds
fi
