# package/network/config/netifd/files/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_toggle_state network "$INTERFACE" up 1
	uci_toggle_state network "$INTERFACE" connect_time $(sed -ne 's![^0-9].*$!!p' /proc/uptime)
	[ -n "$DEVICE" ] && {
		uci_toggle_state network "$INTERFACE" device "$(uci -q get network.$INTERFACE.ifname)"
		uci_toggle_state network "$INTERFACE" ifname "$DEVICE"
	}
}
