# package/network/config/netifd/files/etc/hotplug.d/iface/10-ipv6-static
case "$ACTION" in
	ifup)
		. /lib/functions.sh

		local proto ip6slaac
		config_load network
		config_get proto "$INTERFACE" proto
		config_get_bool ip6slaac "$INTERFACE" ip6slaac

		if [ "$proto" = "static" -a "$ip6slaac" = "1" ]; then
			echo 2 > "/proc/sys/net/ipv6/conf/$DEVICE/accept_ra"

			# Though this should be save here, it is not recommended
			echo 1 > "/proc/sys/net/ipv6/conf/$DEVICE/disable_ipv6"
			echo 0 > "/proc/sys/net/ipv6/conf/$DEVICE/disable_ipv6"
		fi
	;;
esac

