# package/network/config/netifd/files/etc/hotplug.d/iface/10-sysctl
# Skip fake devices (e.g. relayd)
grep -qs "^ *$DEVICE:" /proc/net/dev || exit 0

case "$ACTION" in
	ifup)
		include /lib/network
		scan_interfaces
	;;
esac

