# package/busybox/files/cron
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start () {
	loglevel=$(uci_get "system.@system[0].cronloglevel")
	[ -z "$(ls /etc/crontabs/)" ] && exit 1
	mkdir -p /var/spool/cron
	ln -s /etc/crontabs /var/spool/cron/ 2>/dev/null
	crond -c /etc/crontabs -l ${loglevel:-5}
}

stop() {
	killall -9 crond
}
