# package/busybox/files/cron
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start () {
	loglevel=$(/sbin/uci get "system.@system[0].cronloglevel" 2>/dev/null)
	[ -z "$(ls /etc/crontabs/)" ] && exit 1
	mkdir -p /var/spool/cron
	[ -L /var/spool/cron/crontabs ] || ln -s /etc/crontabs /var/spool/cron/crontabs
	crond -c /etc/crontabs -l ${loglevel:-5}
}

stop() {
	killall -9 crond
}
