# openwrt/package/base-files/default/etc/hotplug.d/net/10-net
# Copyright (C) 2006 OpenWrt.org

include network

addif() {
	scan_interfaces
	setup_interface "$INTERFACE"
	
	# find all vlan configurations for this interface and set them up as well
	for ifc in $interfaces; do
		config_get iftype "$ifc" type
		case "$iftype" in
			bridge) config_get ifs "$ifc" ifnames;;
			*) config_get ifs "$ifc" ifname;;
		esac
		for dev in $ifs; do
			[ "${dev%%\.*}" = "$INTERFACE" -a "$dev" != "$INTERFACE" ] && {
				add_vlan "$dev"
			}
		done
	done
}

case "$ACTION" in
	add|register)
		addif
	;;
esac
