# package/madwifi/files/etc/hotplug.d/net/10-madwifi
if [ "$ACTION" = "add" -o "$ACTION" = "register" ]; then
	case "$INTERFACE" in
		ath*.sta*)
			local BASEIF="${INTERFACE%%\.*}"

			include /lib/network
			scan_interfaces
			local CONFIG="$(find_config "$BASEIF")" 
			[ -n "$CONFIG" ] && setup_interface "$INTERFACE" "$CONFIG"
		;;
	esac
fi
