# package/luci-app-oaf/root/etc/uci-defaults/93_feature_2.0
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set appfilter.feature.update='0'
	set appfilter.feature.format='v2.0'
	commit appfilter
EOF
