# target/linux/mediatek/filogic/base-files/etc/uci-defaults/99-eax17-ledq
#!/bin/sh
# SPDX-License-Identifier: (GPL-2.0 OR MIT)

/etc/init.d/eax17-ledq enable
/etc/init.d/eax17-ledq start
exit 0
