# target/linux/mediatek/filogic/base-files/etc/uci-defaults/99-tpdata-enable
#!/bin/sh
#neccessary for TP-Link Archer AX80 EU version's Wi-Fi to work
/etc/init.d/tpdata enable 2>/dev/null
exit 0
