# target/linux/x86-2.6/base-files/default/etc/inittab
::sysinit:/etc/init.d/rcS
ttyS0::askfirst:/bin/ash --login
tty1::askfirst:/bin/ash --login
