# feeds/luci/themes/luci-theme-openwrt/root/etc/uci-defaults/30_luci-theme-openwrt
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.OpenWrt=/luci-static/openwrt.org
		set luci.main.mediaurlbase=/luci-static/openwrt.org
		commit luci
	EOF
fi

exit 0
