# package/fuse24/files/fuse24.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org
START=15

start() {
	if ls /lib/modules/$(uname -r)/fuse.* > /dev/null 2>&1; then
		[ ! -d /dev/misc ] && umask 022 && mkdir /dev/misc;
		ln -s /dev/fuse /dev/misc/;
	fi
}
