# target/linux/brcm63xx/base-files/etc/init.d/defconfig
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2009 OpenWrt.org
#

START=05

start() {
	. /lib/brcm63xx.sh

	[ ! -d /etc/defconfig/$board_name ] && board_name="generic"

	for f in $( ls /etc/defconfig/$board_name ); do
		if [ ! -e /etc/config/$f ]; then
			cp /etc/defconfig/$board_name/$f /etc/config/
		fi
	done
}
