# target/linux/brcm47xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh
#
# Copyright (C) 2007 OpenWrt.org
#
#

do_fixcrc() {
	mtd fixtrx firmware
}

case `uname -r` in
	3.8* | 3.9* ) do_fixcrc;;
esac

