# target/linux/pistachio/base-files/etc/init.d/set_mac
#!/bin/sh /etc/rc.common

. /lib/functions/mac.sh

START=90
STOP=15

start()
{
    generate_mac $SOFTAP_OTP_REG0
    uci set wireless.ap.macaddr=$MAC
    uci commit wireless
    sync
    /etc/init.d/network restart

    generate_mac $ETH_OTP_REG0
    ifconfig eth0 down
    ifconfig eth0 hw ether $MAC
    ifconfig eth0 up
}

stop()
{
    echo ""
}
