# package/ext-buttons/files/etc/hotplug.d/button/10-buttonchk
#!/bin/sh

log() {
	logger -t "Button Checker " "$@"
}

log "Button Name : $BUTTON     Action : $ACTION"

log "$ACTION $SEEN"

