# target/linux/bcm27xx/base-files/etc/uci-defaults/04-stats-data-dirs.sh
#!/bin/sh

#Create data dir if it does not exist
mkdir -p /data

~
