# target/linux/bcm27xx/base-files/etc/uci-defaults/22-setup-default-ip.sh
#!/bin/sh
uci set network.lan.ipaddr='192.168.8.1'
uci commit network
exit 0
