# target/linux/bcm27xx/base-files/etc/uci-defaults/24-setup-wan.sh
#!/bin/sh
uci set network.wan='interface'
uci set network.wan.ifname='eth1'
uci set network.wan.proto='dhcp'
uci commit network
exit 0
