# package/base-files/files/etc/hotplug.d/usb/10-usb-storage
#!/bin/sh

# Copyright (C) 2006 OpenWrt.org

case "$ACTION" in
	add)
		[ -n "${INTERFACE}" ] &&
		[ "$(expr substr ${INTERFACE} 1 2)" == "8/" ]  && {
			/sbin/usb-storage &
		}
		;;
esac	
		
