# package/base-files/files/etc/hotplug.d/usb/10-usb
#!/bin/sh

# Copyright (C) 2009 OpenWrt.org

case "$ACTION" in
	add)
		# update LEDs
		;;
	remove)
		# update LEDs
		;;
esac	
		
