# target/linux/qualcommbe/ipq95xx/base-files/etc/hotplug.d/ieee80211/11_fix_wifi_mac
#!/bin/sh

[ "$ACTION" == "add" ] || exit 0

[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

set_radio_mac()
{
	local section="$1"
	local band

	config_get band "$section" band
	config_get radio "$section" radio

	case "$board" in
		askey,sbe1v1k)
			WIFIMAC2G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -3)                                    
			WIFIMAC5G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -2)                                    
			WIFIMAC6G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -1)
			
			case "$band" in
				2g)
					uci get wireless.@wifi-iface["$radio"].macaddr || uci set wireless.@wifi-iface["$radio"].macaddr="$WIFIMAC2G"
					;;
				5g)
					uci get wireless.@wifi-iface["$radio"].macaddr || uci set wireless.@wifi-iface["$radio"].macaddr="$WIFIMAC5G"
					;;
				6g)
					uci get wireless.@wifi-iface["$radio"].macaddr || uci set wireless.@wifi-iface["$radio"].macaddr="$WIFIMAC6G"
					;;
			esac
		;;
	esac
}


config_load wireless
config_foreach set_radio_mac wifi-device

uci commit wireless

exit 0

