# package/base-files/files/etc/inittab
::sysinit:/etc/init.d/rcS
tts/0::askfirst:/bin/ash --login
ttyS0::askfirst:/bin/ash --login
