# package/ext-rooter-basic/files/etc/init.d/rooter
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=99

start() {
	/usr/lib/rooter/initialize.sh
}

