# target/linux/ath79/generic/base-files/etc/hotplug.d/firmware/10-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath9k-eeprom-ahb-18100000.wmac.bin")
	case $board in
	avm,fritz1750e|\
	avm,fritz4020|\
	avm,fritz450e|\
	avm,fritzdvbc)
		caldata_extract_reverse "urlader" 0x1541 0x440
		;;
	dlink,dir-842-c1|\
	dlink,dir-842-c2|\
	dlink,dir-842-c3|\
	nec,wf1200cr|\
	nec,wg1200cr|\
	wd,mynet-n600|\
	wd,mynet-n750)
		ath9k_patch_mac $(mtd_get_mac_ascii devdata "wlan24mac")
		;;
	engenius,ecb1200|\
	engenius,ecb1750)
		ath9k_patch_mac $(macaddr_add $(mtd_get_mac_ascii u-boot-env athaddr) 1)
		;;
	enterasys,ws-ap3705i)
		ath9k_patch_mac $(mtd_get_mac_ascii u-boot-env0 RADIOADDR1)
		;;
	extreme-networks,ws-ap3805i)
		ath9k_patch_mac $(mtd_get_mac_ascii cfg1 RADIOADDR1)
		;;
	iodata,wn-ac1167dgr|\
	iodata,wn-ac1600dgr|\
	iodata,wn-ac1600dgr2|\
	iodata,wn-ag300dgr|\
	sitecom,wlr-7100|\
	sitecom,wlr-8100)
		ath9k_patch_mac $(mtd_get_mac_ascii u-boot-env ethaddr)
		;;
	moxa,awk-1137c)
		ath9k_patch_mac $(mtd_get_mac_ascii u-boot-env mac_addr)
		;;
	qihoo,c301)
		ath9k_patch_mac $(mtd_get_mac_ascii devdata "wlan24mac")
		;;
	tplink,deco-s4-v2)
		base_mac=$(mtd_get_mac_encrypted_deco $(find_mtd_part config))
		ath9k_patch_mac $(macaddr_add $base_mac 1)
		;;
	*)
		caldata_die "board $board is not supported yet"
		;;
	esac
	;;
"ath9k-eeprom-pci-0000:00:00.0.bin")
	case $board in
	avm,fritz300e)
		caldata_extract_reverse "urloader" 0x1541 0x440
		;;
	enterasys,ws-ap3705i)
		ath9k_patch_mac $(mtd_get_mac_ascii u-boot-env0 RADIOADDR0)
		;;
	wd,mynet-n600|\
	wd,mynet-n750)
		ath9k_patch_mac $(mtd_get_mac_ascii devdata "wlan5mac")
		;;
	wd,mynet-wifi-rangeextender)
		ath9k_patch_mac $(nvram get wl0_hwaddr)
		;;
	*)
		caldata_die "board $board is not supported yet"
		;;
	esac
	;;
esac
