# target/linux/imx6/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2013 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh

ucidef_set_interface_loopback

case "`cat /proc/device-tree/model`" in
"Gateworks Ventana i.MX6 DualLite/Solo GW51XX" |\
"Gateworks Ventana i.MX6 DualLite/Solo GW52XX" |\
"Gateworks Ventana i.MX6 Dual/Quad GW51XX" |\
"Gateworks Ventana i.MX6 Dual/Quad GW52XX")
	ucidef_set_interface_lan 'eth0'
	;;
"Gateworks Ventana i.MX6 DualLite/Solo GW53XX" |\
"Gateworks Ventana i.MX6 DualLite/Solo GW54XX" |\
"Gateworks Ventana i.MX6 Dual/Quad GW53XX" |\
"Gateworks Ventana i.MX6 Dual/Quad GW54XX" |\
"Gateworks Ventana GW5400-A")
	ucidef_set_interfaces_lan_wan 'eth0' 'eth1'
	;;
"Wandboard i.MX6 Dual Lite Board")
	ucidef_set_interface_wan 'eth0'
	;;
esac

uci commit network

exit 0
