# target/linux/ramips/mt7621/base-files/etc/uci-defaults/04_led_migration
. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
dlink,dir-1960-a1|\
dlink,dir-2640-a1|\
dlink,dir-2660-a1|\
dlink,dir-3040-a1|\
dlink,dir-3060-a1)
	migrate_leds ':net=:wan'
	;;
tplink,archer-a6-v3|\
tplink,archer-c6-v3)
	migrate_leds ':wifi2g$=:wlan-2' ':wifi5g$=:wlan-5'
	;;
esac

remove_devicename_leds

migrations_apply system

exit 0
