# target/linux/ramips/mt7621/base-files/etc/uci-defaults/05_fix-compat-version
#
# Copyright (C) 2020 OpenWrt.org
#

compat_version=$(cat /etc/board.json | jsonfilter -e "@['system']['compat_version']")
[ -n "$compat_version" ] || compat_version="1.0"

uci set system.@system[0].compat_version="$compat_version"
uci commit system

exit 0
