# target/linux/mpc83xx/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2012 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/mpc83xx.sh

ucidef_set_interface_loopback

board=$(mpc83xx_board_name)

case "$board" in
rb333)
	ucidef_set_interfaces_lan_wan "eth0 eth1" "eth2"
	;;

*)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	;;
esac

uci commit network

exit 0
