# target/linux/mvebu/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2014 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/mvebu.sh

board=$(mvebu_board_name)

case "$board" in
armada-xp-mamba)
	ucidef_set_led_netdev "wan" "WAN" "tlc59116:1" "eth1"
	ucidef_set_led_usbdev "usb1" "USB1" "tlc59116:5" "1-1"
	ucidef_set_led_usbdev "usb2" "USB2" "tlc59116:6" "2-2"
	;;

esac

ucidef_commit_leds

exit 0
