# target/linux/mvebu/base-files/etc/uci-defaults/02_network
#!/bin/sh
#
# Copyright (C) 2014-2015 OpenWrt.org
#

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

. /lib/functions/uci-defaults.sh
. /lib/mvebu.sh

ucidef_set_interface_loopback

board=$(mvebu_board_name)

case "$board" in
armada-385-linksys-caiman)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-385-linksys-cobra)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-385-linksys-shelby)
	ucidef_set_interfaces_lan_wan "eth1" "eth0"
	;;
armada-xp-linksys-mamba)
	ucidef_set_interfaces_lan_wan "eth0" "eth1"
	ucidef_add_switch "switch0" "1" "1"
	ucidef_add_switch_vlan "switch0" "1" "0 1 2 3 5"
	ucidef_add_switch_vlan "switch0" "2" "4 6"
	;;
armada-385-db-ap)
	ucidef_set_interfaces_lan_wan "eth0 eth1" "eth2"
	;;
armada-xp-gp)
	ucidef_set_interface_lan "eth0 eth1 eth2 eth3"
	;;
armada-388-clearfog)
	ucidef_set_interfaces_lan_wan "eth0 eth1" "eth2"
	;;
*)
	ucidef_set_interface_lan "eth0"
	;;
esac

uci commit network

exit 0
