# target/linux/mvebu/base-files/etc/uci-defaults/03_wireless
#!/bin/sh
#
# Copyright (C) 2015 OpenWrt.org
#

[ ! -e /etc/config/wireless ] && exit 0

. /lib/mvebu.sh

board=$(mvebu_board_name)

case "$board" in
armada-385-linksys-caiman|armada-385-linksys-cobra|armada-385-linksys-shelby|armada-xp-linksys-mamba)
	SKU=$(strings /dev/mtd3|sed -ne 's/^cert_region=//p')
	case "$SKU" in
		AP)
			REGD=CN
		;;
		AU)
			REGD=AU
		;;
		CA)
			REGD=CA
		;;
		EU)
			REGD=DE
		;;
		US)
			REGD=US
		;;
	esac
	uci get wireless.radio0.country || uci set wireless.radio0.country=$REGD
	uci get wireless.radio1.country || uci set wireless.radio1.country=$REGD
	;;
esac

uci commit wireless

exit 0
