# root/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/sbin/halt
::respawn:/usr/sbin/dnsmasq -d -l /tmp/udhcpd.leases -s lan
