# package/kernel/8563-watchdog/files/8563-watchdog.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2011 OpenWrt.org
#

START=11
STOP=11

USE_PROCD=1
NAME=8563-watchdog
PROG=/sbin/8563-watchdog

start_service() {
	procd_open_instance
	procd_set_param command "${PROG}" "${gpio}"
	procd_set_param respawn
	procd_close_instance
}
