# feeds/luci/applications/luci-app-nextdns/root/etc/uci-defaults/60_luci-nextdns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@nextdns[-1]
	add ucitrack nextdns
	set ucitrack.@nextdns[-1].init=nextdns
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
