# xnd/wifi-ax/mac80211/files/etc/init.d/ath11k
#!/bin/sh /etc/rc.common

START=11

boot() {
	sleep 5
	modprobe ath11k_ahb
}
