# target/linux/lantiq/base-files/etc/init.d/esi
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=19
start() {
	esi $(printf '%X' $((1+0x$(tr -d : </sys/class/net/eth0/address)))) 2>/dev/null || :
}
