# package/independent/lte-modem-xwrt-csac/files/lte-modem.init
#!/bin/sh /etc/rc.common
# (C) 2012 Daniel Golle, Allnet GmbH <dgolle@allnet.de>

START=19
STOP=90

boot() {
	echo "0" >/sys/class/gpio/lte_modem_enable/value
	sleep 1
	echo "1" >/sys/class/gpio/lte_modem_enable/value
}
