# files/etc/hotplug.d/iface/99-odhcpd-qmi
#!/bin/sh
[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "qmi" ] || exit 0
sleep 5
/etc/init.d/odhcpd restart
