# target/linux/mediatek/filogic/base-files/etc/hotplug.d/firmware/11-mt76-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"mediatek/mt7916_eeprom.bin")
	case "$board" in
	acer,predator-w6)
		caldata_extract_mmc "factory" 0xA0000 0x1000
		;;
	esac
	;;
"mediatek/mt7981_eeprom_mt7976_dbdc.bin")
	case "$board" in
	cmcc,rax3000m)
		case "$(cmdline_get_var root)" in
		/dev/mmc*)
			caldata_extract_mmc "factory" 0x0 0x1000
			;;
		esac
		;;
	ubnt,unifi-6-plus)
		caldata_extract_mmc "factory" 0x0 0x1000
		;;
	esac
	;;
"mediatek/mt7986_eeprom_mt7975_dual.bin")
	case "$board" in
	mercusys,mr90x-v1)
		ln -sf /tmp/tp_data/MT7986_EEPROM.bin \
			/lib/firmware/$FIRMWARE
		;;
	esac
	;;
"mediatek/mt7986_eeprom_mt7976.bin")
	case "$board" in
	acer,predator-w6)
		caldata_extract_mmc "factory" 0x0 0x1000
		;;
	esac
	;;
"mediatek/mt7986_eeprom_mt7976_dbdc.bin")
	case "$board" in
	asus,tuf-ax4200|\
	asus,tuf-ax6000)
		CI_UBIPART="UBI_DEV"
		caldata_extract_ubi "Factory" 0x0 0x1000
		;;
	esac
	;;
"mediatek/mt7986_eeprom_mt7976_dual.bin")
	case "$board" in
	glinet,gl-mt6000)
		caldata_extract_mmc "factory" 0x0 0x1000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
