# target/linux/qualcommbe/ipq95xx/base-files/etc/uci-defaults/03_wireless
#!/bin/sh

[ ! -e /etc/config/wireless ] && exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
askey,sbe1v1k)
	WIFIMAC2G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -3)
	WIFIMAC5G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -2)
	WIFIMAC6G=$(macaddr_add $(mmc_get_mac_binary 0:ART 0) -1)

	uci get wireless.@wifi-iface[0].macaddr || uci set wireless.@wifi-iface[0].macaddr=$WIFIMAC2G
	uci get wireless.@wifi-iface[1].macaddr || uci set wireless.@wifi-iface[1].macaddr=$WIFIMAC5G
	uci get wireless.@wifi-iface[2].macaddr || uci set wireless.@wifi-iface[2].macaddr=$WIFIMAC6G
	;;
esac

uci commit wireless

exit 0
