# package/rooter/0optionalapps/luci-app-hotspot/files/etc/init.d/zhot
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=99

log() {
	wifilog "Z-Hotspot " "$@"
}

start()
{
	uci set travelmate.global.lost="0"
	uci set travelmate.global.delay="30"
	uci commit travelmate
	/usr/lib/hotspot/travelmate.sh &
}
