# openwrt/package/base-files/default/etc/init.d/S40network
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

start() {
	setup_switch() { return 0; }

	include /lib/network
	setup_switch
	/sbin/wifi
}

