# xnd/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-opennds
#!/bin/sh

uci delete opennds.@opennds[-1]
