# xnd/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-wifi
#!/bin/sh

. /lib/functions.sh

iface_del() {
        uci delete wireless.$1
}

config_load wireless
config_foreach iface_del wifi-iface
