# xnd/ucentral/ucentral-schema/files/etc/uci-defaults/99-ucentral-umdns
#!/bin/sh

uci delete umdns.@umdns[0].network
