# openwrt/package/lighttpd/ipkg/lighttpd.conffiles
/etc/lighttpd.conf
/etc/default/lighttpd
/etc/init.d/lighttpd
