# package/system/procd/files/procd.init
#!/bin/sh /etc/rc.common

START=11

start_stop() {
	start-stop-daemon $1 -b -m -p /var/run/procd.pid -x /sbin/procd
}

start() {
	start_stop -S
}

reload() {
	return
}

stop() {
	start_stop -K
}
