# package/base-files/files/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_set_state network "$INTERFACE" up 1 
}
