# package/base-files/files/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_set_state network "$INTERFACE" up 1 
	[ -n "$DEVICE" ] && uci_set_state network "$INTERFACE" ifname "$DEVICE"
}
