# package/base-files/files/etc/hotplug.d/iface/00-netstate
[ ifup = "$ACTION" ] && {
	uci set "/var/state/network.$INTERFACE.up=1"
}
