# package/om2p-watchdog/files/om2p-watchdog.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2011 OpenWrt.org
#

START=11

boot() {
	. /lib/ar71xx.sh

	local board=$(ar71xx_board_name)

	[ "$board" = "om2p" ] && /sbin/om2p-watchdog &
}
