# target/linux/ath79/tiny/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
#!/bin/ash

[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

board=$(board_name)

case "$board" in
	dlink,dir-615-e4)
		[ "$PHYNBR" -eq 0 ] && \
			mtd_get_mac_ascii "nvram" "lan_mac" > /sys${DEVPATH}/macaddress
		;;
esac
